# kb-content-output/articles/incident_response/eradication-procedures.json
{
  "title": "Eradication Procedures",
  "content": "## Purpose\n\nThis runbook provides detailed procedures for eradicating threat actor presence from compromised systems. Eradication ensures all persistence mechanisms, malware artifacts, unauthorized accounts, and attacker footholds are completely removed before recovery begins.\n\n## Scope\n\nApplies to all confirmed incidents where a threat actor has established persistence, deployed malware, or made unauthorized modifications to systems. Must be executed after containment is confirmed effective and forensic evidence has been preserved.\n\n## Prerequisites\n\n- Incident has been contained and containment verified (no active C2, no ongoing lateral movement)\n- Forensic evidence has been collected and preserved per the Evidence Preservation SOP\n- Full scope of compromise has been determined (all affected systems, accounts, and data identified)\n- EDR, SIEM, and forensic tooling access available\n- Incident Commander has approved eradication plan\n\n## Procedure\n\n### Step 1: Persistence Mechanism Identification\n\nBefore removing anything, conduct a comprehensive audit of all persistence mechanisms the threat actor may have established. Check each of the following on every affected system:\n\n1. **Windows persistence locations:**\n   - Registry Run/RunOnce keys: `HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run`, `HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run`\n   - Scheduled Tasks: `schtasks /query /fo LIST /v` — look for unfamiliar tasks, especially those running as SYSTEM\n   - Services: Review recently created or modified services with `sc query` and check for services with suspicious binary paths\n   - WMI Event Subscriptions: `Get-WMIObject -Namespace root\\subscription -Class __EventFilter` and related classes\n   - Startup folder items: `%APPDATA%\\Microsoft\\Windows\\Start Menu\\Programs\\Startup`\n   - DLL search order hijacking: Verify DLLs in application directories match expected hashes\n   - COM object hijacking: Check `HKCU\\Software\\Classes\\CLSID` for unexpected entries\n   - Print monitor DLLs: `HKLM\\SYSTEM\\CurrentControlSet\\Control\\Print\\Monitors`\n   - AppInit_DLLs: `HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Windows\\AppInit_DLLs`\n\n2. **Linux persistence locations:**\n   - Cron jobs: `/etc/crontab`, `/var/spool/cron/`, `/etc/cron.d/`, user crontabs\n   - Systemd services and timers: `/etc/systemd/system/`, `~/.config/systemd/user/`\n   - SSH authorized_keys: `~/.ssh/authorized_keys` for all user accounts including root\n   - Shell profile scripts: `.bashrc`, `.bash_profile`, `.profile`, `/etc/profile.d/`\n   - LD_PRELOAD: Check `/etc/ld.so.preload` and environment variables\n   - Kernel modules: `lsmod` and check `/etc/modules-load.d/`\n   - At jobs: `atq` and `/var/spool/at/`\n\n3. **Cloud and identity persistence:**\n   - Rogue user accounts or service principals in Active Directory, Azure AD/Entra ID, AWS IAM\n   - OAuth application consents and enterprise app registrations\n   - Federated identity trust relationships (Golden SAML attack vector)\n   - API keys, access keys, and long-lived tokens\n   - Email forwarding rules and mailbox delegation\n   - Conditional access policy modifications\n\n4. **Web-based persistence:**\n   - Web shells in web server document roots and upload directories\n   - Modified web application configuration files\n   - Rogue reverse proxies or SOCKS tunnels\n   - Injected JavaScript in web application source files\n\n### Step 2: Malware Removal\n\n1. Using the inventory from Step 1, remove all identified malware binaries, scripts, and configuration files from affected systems.\n2. For each file removed, record the full file path, SHA-256 hash, file size, creation/modification timestamps, and the user context under which it was running.\n3. Quarantine samples rather than deleting them — store in a password-protected ZIP in the forensics evidence repository.\n4. Remove all persistence mechanisms identified in Step 1: delete scheduled tasks, registry keys, services, cron jobs, startup items, and WMI subscriptions.\n5. For web shells, remove the file and compare the web root against a known-good baseline to identify any other modified files.\n6. Check for rootkits using specialized tools (GMER, rkhunter, chkrootkit). If a kernel-level rootkit is detected, the system must be rebuilt from scratch.\n\n### Step 3: Credential Rotation\n\nAssume all credentials on compromised systems are compromised. Execute the following:\n\n1. **User accounts**: Force password reset for all users whose workstations or accounts were affected. Require MFA re-enrollment.\n2. **Service accounts**: Rotate passwords for all service accounts that were present on compromised systems. Coordinate with application teams to update dependent configurations.\n3. **Kerberos**: If a domain controller was compromised, perform a double KRBTGT password reset (with 12-hour interval between resets) to invalidate all Golden Tickets.\n4. **Certificates**: If private keys were stored on compromised systems, revoke the associated certificates and issue new ones.\n5. **API keys and tokens**: Rotate all API keys, OAuth client secrets, and access tokens associated with compromised systems or accounts.\n6. **Database credentials**: Change database passwords and connection strings for any database accessed from compromised systems.\n7. **SSH keys**: Regenerate SSH key pairs and remove old public keys from authorized_keys files across the environment.\n\n### Step 4: Vulnerability Remediation\n\n1. Identify the initial access vector — the vulnerability, misconfiguration, or social engineering technique that enabled the compromise.\n2. Patch the exploited vulnerability across all systems, not just the affected ones.\n3. If the initial access was via phishing, update email gateway rules, conduct user awareness training, and review email authentication (SPF, DKIM, DMARC).\n4. If the initial access was via an exposed service, review firewall rules and reduce the external attack surface.\n5. If the initial access was via a supply chain compromise, evaluate the vendor relationship and implement additional monitoring.\n\n### Step 5: Eradication Verification\n\n1. Run full antivirus/EDR scans with updated signatures on all affected systems.\n2. Compare current system state against known-good baselines (file integrity monitoring, configuration baselines).\n3. Review network traffic from affected systems for any residual C2 activity.\n4. Verify that all credential rotations have been completed and old credentials are no longer functional.\n5. Conduct a targeted threat hunt across the environment using the IOCs and TTPs from this incident to confirm no additional compromised systems exist.\n6. Have a second analyst independently verify eradication before proceeding to recovery.\n\n## Escalation Criteria\n\n- Kernel-level rootkit or firmware implant detected: Escalate to specialized forensics team; system must be physically replaced or rebuilt from bare metal\n- Active Directory compromise with KRBTGT exposure: Escalate to Identity team lead; forest-wide credential rotation may be required\n- Supply chain or vendor-originated compromise: Escalate to CISO for vendor communication and potential legal action\n- Eradication actions will require extended downtime for critical systems: Escalate to business leadership for approval\n\n## References\n\n- NIST SP 800-61 Rev. 2, Section 3.3.4 — Eradication\n- NIST SP 800-83 — Guide to Malware Incident Prevention and Handling\n- MITRE ATT&CK — Persistence Tactic (TA0003)\n- Microsoft Advanced Threat Analytics — Golden Ticket Detection\n- CIS Controls v8 — Control 10: Malware Defenses",
  "content_type": "runbook",
  "category": "incident_response",
  "subcategory": "eradication",
  "tags": ["eradication", "persistence", "malware-removal", "credential-rotation", "KRBTGT", "web-shell", "rootkit", "incident-response"],
  "severity_filter": ["high", "critical"],
  "incident_types": ["malware", "ransomware", "unauthorized-access", "advanced-persistent-threat", "web-shell"],
  "ioc_types": ["hash", "ip", "domain", "url"],
  "mitre_techniques": ["T1053.005", "T1547.001", "T1543.003", "T1546.003", "T1136.001", "T1098", "T1505.003", "T1556.001", "T1558.001"],
  "compliance_frameworks": ["NIST", "SOC2", "ISO27001", "PCI-DSS"],
  "priority": 85
}
