# kb-content-output/articles/incident_response/ir-lifecycle-overview.json
{
  "title": "Incident Response Lifecycle Overview",
  "content": "## Purpose\n\nThis document defines the standard Incident Response (IR) lifecycle used by the SOC team, aligned with NIST SP 800-61 Rev. 2. It establishes the phases of incident handling, defines when to activate the formal IR plan, and provides decision criteria for escalation between phases.\n\n## Scope\n\nApplies to all security incidents detected through automated alerting, threat intelligence, user reports, or third-party notifications. All SOC analysts (Tier 1 through Tier 3), IR leads, and incident commanders must follow this lifecycle.\n\n## Prerequisites\n\n- Access to the T1 Agentics SOAR platform and case management system\n- Familiarity with the organization's asset inventory and network topology\n- Current incident classification matrix (see Severity Classification Matrix)\n- Contact information for escalation chain and on-call personnel\n\n## Procedure\n\n### Phase 1: Preparation\n\n1. Maintain and review IR playbooks quarterly. Ensure all runbooks reflect current infrastructure, tools, and team structure.\n2. Verify that detection rules, SIEM correlation logic, and alerting thresholds are tuned and generating actionable alerts.\n3. Confirm all analysts have valid credentials for forensic tools, EDR consoles, cloud admin panels, and communication channels.\n4. Conduct tabletop exercises at least twice per year to validate response readiness.\n5. Maintain a pre-staged IR toolkit including forensic imaging software, write blockers, chain-of-custody forms, and network isolation scripts.\n6. Review and update the asset criticality matrix to ensure business impact assessments reflect current priorities.\n\n### Phase 2: Detection and Analysis\n\n1. Monitor all alert sources: SIEM, EDR, NDR, email gateway, cloud security posture management (CSPM), identity threat detection, and user-reported incidents.\n2. When an alert fires, perform initial triage within 15 minutes for high/critical alerts and 60 minutes for medium/low alerts.\n3. Correlate the alert with related telemetry — look for the same source IP, user account, hostname, or file hash across multiple data sources.\n4. Determine if the event is a true positive, false positive, or benign true positive. Document your reasoning.\n5. If confirmed as a true positive, classify severity using the Incident Severity Classification Matrix (P1-P4).\n6. For P1/P2 incidents, immediately activate the IR plan and notify the Incident Commander.\n7. Create a formal incident case in the SOAR platform, linking all related alerts, observables, and analyst notes.\n\n### Phase 3: Containment\n\n1. Select a containment strategy appropriate to the threat type and business impact. Options include:\n   - **Short-term containment**: Network isolation of affected host, disabling compromised account, blocking malicious IP/domain at firewall.\n   - **Long-term containment**: Applying temporary security controls, deploying additional monitoring on affected segments, implementing network segmentation changes.\n2. Preserve forensic evidence BEFORE taking containment actions that may alter system state. Take memory dumps, disk images, and log snapshots.\n3. Document all containment actions with timestamps, the analyst who performed them, and the business justification.\n4. Validate containment effectiveness — confirm the threat actor no longer has active access by monitoring for continued C2 beaconing, lateral movement, or data exfiltration.\n5. Communicate containment status to stakeholders per the Incident Communication Plan.\n\n### Phase 4: Eradication\n\n1. Identify all persistence mechanisms: scheduled tasks, registry run keys, cron jobs, web shells, implanted SSH keys, rogue user accounts, or modified startup scripts.\n2. Remove all malware, backdoors, and unauthorized changes from affected systems.\n3. Rotate all credentials that may have been exposed — this includes service accounts, API keys, database passwords, and user passwords.\n4. Patch the vulnerability or close the attack vector that enabled initial access.\n5. Verify eradication by scanning with updated signatures and reviewing system baselines.\n\n### Phase 5: Recovery\n\n1. Restore affected systems from known-good backups or rebuild from golden images.\n2. Reintroduce systems to the production environment in a staged manner — do not restore all systems simultaneously.\n3. Implement enhanced monitoring on recovered systems for a minimum of 30 days.\n4. Validate that business services are functioning correctly and that no indicators of re-compromise appear.\n5. Obtain sign-off from the system owner and IR lead before declaring recovery complete.\n\n### Phase 6: Post-Incident Activity\n\n1. Conduct a post-incident review (lessons learned) within 5 business days of incident closure.\n2. Reconstruct a detailed timeline from initial compromise through detection, containment, eradication, and recovery.\n3. Identify gaps in detection, response, or communication and create action items with owners and deadlines.\n4. Update playbooks, detection rules, and training materials based on findings.\n5. Produce a final incident report and distribute to appropriate stakeholders.\n\n## IR Plan Activation Criteria\n\nActivate the formal IR plan when ANY of the following conditions are met:\n\n- A P1 or P2 incident is confirmed\n- Active data exfiltration is detected\n- Ransomware encryption is observed on any production system\n- A threat actor has established persistence in the environment\n- A zero-day vulnerability is being actively exploited\n- Law enforcement or a regulatory body contacts the organization about a compromise\n- A third-party vendor reports a breach affecting shared data or systems\n\n## Escalation Criteria\n\n| Condition | Action |\n|-----------|--------|\n| Alert cannot be triaged within SLA | Escalate to Tier 2 |\n| Confirmed compromise of critical asset | Escalate to Incident Commander |\n| Potential data breach involving PII/PHI | Notify Legal and Privacy Officer |\n| Incident requires law enforcement | Notify CISO and General Counsel |\n| Containment requires business service interruption | Obtain CAB/management approval |\n\n## References\n\n- NIST SP 800-61 Rev. 2 — Computer Security Incident Handling Guide\n- NIST Cybersecurity Framework (CSF) v2.0 — Respond Function\n- SANS Incident Handler's Handbook\n- MITRE ATT&CK Framework — Enterprise Matrix\n- Organization-specific IR Plan (controlled document)",
  "content_type": "sop",
  "category": "incident_response",
  "subcategory": "lifecycle",
  "tags": ["incident-response", "NIST-800-61", "IR-plan", "lifecycle", "preparation", "detection", "containment", "eradication", "recovery"],
  "severity_filter": ["low", "medium", "high", "critical"],
  "incident_types": ["malware", "phishing", "data-breach", "unauthorized-access", "ransomware", "insider-threat", "denial-of-service"],
  "ioc_types": ["ip", "domain", "hash", "url", "email"],
  "mitre_techniques": ["T1059", "T1053", "T1078", "T1486", "T1041", "T1071"],
  "compliance_frameworks": ["NIST", "SOC2", "ISO27001", "PCI-DSS", "HIPAA", "GDPR"],
  "priority": 95
}
