# kb-content-output/articles/threat_detection/lateral-movement-detection.json
{
  "title": "Lateral Movement Detection Guide",
  "content": "## Purpose\n\nThis SOP defines the detection strategies and investigation procedures for identifying lateral movement within the enterprise network. Lateral movement is a post-compromise technique where adversaries move through the environment to reach target systems, escalate privileges, and access sensitive data. Early detection of lateral movement significantly reduces dwell time and limits blast radius.\n\n## Scope\n\nCovers detection of lateral movement via common protocols and tools including SMB/PsExec, WMI, RDP, SSH, WinRM, DCOM, and pass-the-hash/pass-the-ticket attacks across Windows and Linux environments.\n\n## Prerequisites\n\n- Windows Security Event Logs and Sysmon deployed across endpoints\n- Network traffic monitoring (NetFlow, full PCAP, or NDR solution)\n- Authentication logs from domain controllers centralized in SIEM\n- EDR telemetry with process creation and network connection visibility\n- Baseline of legitimate administrative tool usage and remote access patterns\n\n## Procedure\n\n### Step 1: Monitor Authentication Patterns\n\n1. **Windows Logon Events**: Monitor Event ID 4624 with special attention to:\n   - Logon Type 3 (Network): Standard for SMB, PsExec, WMI remote access\n   - Logon Type 10 (RemoteInteractive): RDP sessions\n   - Logon Type 7 (Unlock): Can indicate RDP session reconnection\n2. **Watch for anomalous authentication patterns**:\n   - Single account authenticating to many hosts in rapid succession (> 5 hosts in 10 minutes)\n   - Authentication from hosts that are not management workstations to servers\n   - Service account authenticating interactively (should only have Logon Type 5)\n   - Local administrator account used across multiple machines (lateral movement via shared local admin passwords)\n3. **Kerberos indicators**:\n   - Event ID 4768 (TGT Request): Normal, but watch for encryption type downgrades (RC4 instead of AES = potential overpass-the-hash)\n   - Event ID 4769 (Service Ticket Request): Excessive requests from a single source\n   - Event ID 4771 (Kerberos Pre-Auth Failure): May indicate password spraying preceding lateral movement\n\n### Step 2: Detect PsExec and SMB-Based Movement\n\n1. **PsExec indicators**:\n   - Service creation event (Event ID 7045) with service name matching `PSEXESVC` or random 8-character names\n   - Named pipe creation: `\\\\*\\PIPE\\PSEXESVC` or `\\\\*\\PIPE\\svcctl`\n   - Sysmon Event ID 17/18 (Pipe Created/Connected) for non-standard pipes\n   - Process tree: `services.exe` > `PSEXESVC.exe` > child process (cmd.exe, powershell.exe)\n2. **SMB file writes to admin shares**:\n   - Event ID 5145 (Network Share Access) targeting `ADMIN$`, `C$`, or `IPC$`\n   - Unexpected binary files written to `\\Windows\\Temp` via SMB\n3. **Impacket-based variants**: Tools like `smbexec`, `atexec`, and `wmiexec` from the Impacket toolkit produce slightly different artifacts:\n   - `smbexec`: Creates a service with a random name, uses `cmd.exe /Q /c` command pattern\n   - `atexec`: Creates a scheduled task (Event ID 4698) instead of a service\n\n### Step 3: Detect WMI-Based Movement\n\n1. **WMI process creation**:\n   - Parent process `WmiPrvSE.exe` spawning `cmd.exe`, `powershell.exe`, or other unexpected children\n   - Event ID 4688 showing `WmiPrvSE.exe` as creator process\n2. **WMI event subscription persistence** (also T1546.003):\n   - Sysmon Event ID 19/20/21 (WmiEventFilter, WmiEventConsumer, WmiEventConsumerToFilter)\n3. **Network indicators**:\n   - DCOM traffic on TCP 135 followed by high ephemeral port connections\n   - WinRM traffic on TCP 5985 (HTTP) or 5986 (HTTPS)\n\n### Step 4: Detect RDP Lateral Movement\n\n1. **Source host indicators**:\n   - Event ID 1024 (RDP ClientActiveX: Attempting to connect) in `Microsoft-Windows-TerminalServices-RDPClient/Operational`\n   - `mstsc.exe` execution (Sysmon Event ID 1)\n2. **Destination host indicators**:\n   - Event ID 4624 Logon Type 10 on the target\n   - Event ID 21 (Session Logon Succeeded) in `Microsoft-Windows-TerminalServices-LocalSessionManager/Operational`\n3. **Suspicious RDP behaviors**:\n   - RDP connections originating from servers (non-admin workstations)\n   - RDP to hosts that have never received RDP connections before\n   - RDP tunneled over SSH (detect by SSH sessions with subsequent RDP traffic)\n   - Multiple RDP sessions from a single source in rapid succession (hopping)\n\n### Step 5: Detect Pass-the-Hash and Pass-the-Ticket\n\n1. **Pass-the-Hash (PtH)** (T1550.002):\n   - Event ID 4624 with Logon Type 9 (NewCredentials) from non-admin tools\n   - NTLM authentication (Event ID 4624 with Authentication Package = NTLM) in an environment that should be using Kerberos\n   - Sysmon Event ID 10 (Process Access) targeting `lsass.exe` with access flags `0x1010` or `0x1038`\n2. **Pass-the-Ticket (PtT)** (T1550.003):\n   - Kerberos TGS requests (Event ID 4769) without a preceding TGT request (Event ID 4768) from the same host\n   - Ticket encryption type anomalies\n   - Tools: Mimikatz `sekurlsa::pth`, Rubeus `ptt`\n3. **Overpass-the-Hash**: NTLM hash used to request Kerberos TGT. Look for Event ID 4768 with RC4 encryption (0x17) when AES should be standard.\n\n### Step 6: Detect Linux/SSH Lateral Movement\n\n1. **SSH-based movement**:\n   - Auth.log entries for SSH key-based authentication from unexpected sources\n   - SSH connections from non-jump-box systems\n   - SSH tunneling: `-L`, `-R`, `-D` flags in command line or connections to localhost ports that redirect to remote hosts\n2. **Linux-specific indicators**:\n   - Use of `scp`, `rsync`, or `sftp` to transfer tools to remote hosts\n   - Remote command execution via `ssh user@host 'command'`\n   - Cron job creation on remote hosts via SSH\n\n## Detection Rules Summary\n\n| Technique | Key Data Source | Detection Logic |\n|---|---|---|\n| PsExec | Event 7045, Sysmon 17/18 | New service + named pipe creation |\n| WMI | Event 4688, Sysmon 1 | WmiPrvSE.exe spawning shell |\n| RDP | Event 4624 Type 10 | RDP from non-admin workstation |\n| PtH | Event 4624 NTLM | NTLM logon in Kerberos environment |\n| PtT | Event 4769 | TGS without preceding TGT |\n| SSH Tunnel | Auth.log, NetFlow | SSH with port forwarding flags |\n\n## Escalation Criteria\n\n- Confirmed lateral movement from any compromised host: Immediate P1 escalation\n- Pass-the-hash or pass-the-ticket activity detected: Engage incident response and Active Directory security team\n- Multiple hosts communicating with a single suspicious internal host: Possible attacker pivot point\n- Any lateral movement touching domain controllers, certificate authorities, or backup infrastructure\n\n## References\n\n- MITRE ATT&CK: Lateral Movement (TA0008)\n- MITRE ATT&CK: T1021 (Remote Services), T1550 (Use Alternate Authentication Material)\n- Microsoft Advanced Threat Analytics documentation\n- SANS: Detecting Lateral Movement through Tracking Event Logs",
  "content_type": "sop",
  "category": "threat_detection",
  "subcategory": "lateral_movement",
  "tags": ["lateral-movement", "psexec", "wmi", "rdp", "ssh", "pass-the-hash", "pass-the-ticket", "mimikatz", "kerberos"],
  "severity_filter": ["high", "critical"],
  "incident_types": ["lateral_movement", "credential_theft", "active_intrusion", "apt"],
  "ioc_types": ["ip", "hash"],
  "mitre_techniques": ["T1021", "T1021.001", "T1021.002", "T1021.003", "T1021.004", "T1021.006", "T1047", "T1550.002", "T1550.003", "T1570", "T1563.002"],
  "compliance_frameworks": ["NIST", "SOC2", "ISO27001", "PCI-DSS"],
  "priority": 90
}
