# kb-content-output/articles/threat_detection/persistence-mechanism-detection.json
{
  "title": "Persistence Mechanism Detection Guide",
  "content": "## Purpose\n\nThis runbook provides step-by-step procedures for detecting persistence mechanisms established by adversaries on compromised endpoints. Persistence ensures the adversary maintains access to the environment across system reboots, credential changes, and other interruptions. Identifying and removing all persistence mechanisms is critical for complete incident remediation.\n\n## Scope\n\nCovers detection of persistence via registry run keys, scheduled tasks, WMI event subscriptions, startup folders, Windows services, DLL hijacking, boot/logon scripts, and Linux equivalents (cron, systemd, init scripts, SSH authorized keys).\n\n## Prerequisites\n\n- Sysmon deployed with registry and WMI event monitoring enabled\n- Windows Event Log audit policies for Object Access and Process Creation\n- Autoruns by Sysinternals available for manual analysis\n- EDR with persistence monitoring capabilities\n- Linux auditd configured for file integrity monitoring on critical paths\n- Baseline of legitimate startup items, services, and scheduled tasks\n\n## Runbook Steps\n\n### Step 1: Registry Run Key Monitoring (T1547.001)\n\nRegistry run keys execute programs when a user logs on or the system starts.\n\n**Key registry locations to monitor**:\n```\n# User-level (per-user persistence)\nHKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\nHKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\RunOnce\nHKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\RunServices\nHKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\RunServicesOnce\nHKCU\\Software\\Microsoft\\Windows NT\\CurrentVersion\\Windows\\Load\nHKCU\\Software\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon\\Shell\nHKCU\\Environment\\UserInitMprLogonScript\n\n# Machine-level (all-user persistence, requires admin)\nHKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\nHKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\RunOnce\nHKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Policies\\Explorer\\Run\nHKLM\\Software\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon\\Shell\nHKLM\\Software\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon\\Userinit\nHKLM\\SYSTEM\\CurrentControlSet\\Control\\Session Manager\\BootExecute\n```\n\n**Detection**:\n1. Sysmon Event ID 13 (Registry Value Set): Alert on any modification to the above keys\n2. Compare new registry values against the baseline of known legitimate entries\n3. Investigate binaries referenced by new entries: Check file hash, digital signature, file creation timestamp, and path\n4. Red flags: Entries pointing to `%TEMP%`, `%APPDATA%`, `C:\\ProgramData`, or user profile directories; unsigned executables; obfuscated paths using environment variables\n\n### Step 2: Scheduled Task Detection (T1053.005)\n\n1. **Event-based detection**:\n   - Event ID 4698 (Scheduled Task Created): Alert on all new task creation\n   - Event ID 4702 (Scheduled Task Updated): Alert on modifications to existing tasks\n   - Sysmon Event ID 1: `schtasks.exe` execution with `/create` parameter\n2. **Suspicious task characteristics**:\n   - Tasks running as SYSTEM or domain admin accounts created by non-admin users\n   - Tasks executing binaries from temp directories, recycler, or user-writable paths\n   - Hidden tasks (tasks with the Hidden XML attribute set to true)\n   - Tasks triggered at logon, idle, or at frequent intervals (< 5 minutes)\n   - XML task definitions with encoded commands or download cradles\n3. **Enumeration command**:\n   - `schtasks /query /fo CSV /v` - Export all tasks for review\n   - Check `C:\\Windows\\System32\\Tasks\\` directory for task XML files\n   - Compare against known-good task baseline\n\n### Step 3: WMI Event Subscription Detection (T1546.003)\n\nWMI persistence is stealthy because it does not create files on disk in some implementations.\n\n1. **Sysmon WMI events**:\n   - Event ID 19 (WmiEventFilter activity detected): New event filter created\n   - Event ID 20 (WmiEventConsumer activity detected): New event consumer created\n   - Event ID 21 (WmiEventConsumerToFilter activity detected): Consumer bound to filter\n2. **Manual enumeration**:\n   ```powershell\n   Get-WMIObject -Namespace root\\Subscription -Class __EventFilter\n   Get-WMIObject -Namespace root\\Subscription -Class __EventConsumer\n   Get-WMIObject -Namespace root\\Subscription -Class __FilterToConsumerBinding\n   ```\n3. **Suspicious indicators**:\n   - CommandLineEventConsumer executing scripts or binaries\n   - ActiveScriptEventConsumer with embedded VBScript or JScript\n   - Event filters triggered on common events: process start, user logon, timer intervals\n   - Any WMI subscription is suspicious in most environments - very few legitimate applications use this mechanism\n\n### Step 4: Startup Folder Monitoring (T1547.001)\n\n1. **Monitored paths**:\n   - `C:\\Users\\<username>\\AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\\` (per-user)\n   - `C:\\ProgramData\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\\` (all users)\n2. **Detection**:\n   - Sysmon Event ID 11 (File Created): Alert on new files in startup folders\n   - Check file types: `.exe`, `.bat`, `.cmd`, `.vbs`, `.js`, `.lnk` (shortcut files pointing to malicious binaries)\n   - LNK files are particularly common: parse the target path of .lnk files to identify the actual executable\n\n### Step 5: Windows Service Persistence (T1543.003)\n\n1. **Service creation**:\n   - Event ID 7045 (Service Installed): Alert on all new services\n   - Event ID 4697 (A service was installed in the system)\n   - Compare against known-good service baseline\n2. **Suspicious service characteristics**:\n   - Service binary path in unusual locations (not `C:\\Windows\\System32\\` or `C:\\Program Files\\`)\n   - Unsigned service binaries\n   - Service names mimicking legitimate services with slight misspellings (e.g., `svch0st` vs `svchost`)\n   - Services set to automatic start with short random names\n   - Service binary path containing `cmd.exe /c` or `powershell.exe` commands\n3. **Service DLL persistence**:\n   - Services using `svchost.exe -k` with a custom service group\n   - Check `HKLM\\SYSTEM\\CurrentControlSet\\Services\\<name>\\Parameters\\ServiceDll` for unexpected DLLs\n\n### Step 6: DLL Search Order Hijacking (T1574.001)\n\n1. **Detection approach**:\n   - Monitor for unsigned DLLs loaded by signed Microsoft binaries (Sysmon Event ID 7)\n   - Alert on DLLs loaded from the current working directory when a system directory version exists\n   - Maintain a list of commonly hijacked DLLs and their expected paths\n2. **Common targets**: `version.dll`, `dbghelp.dll`, `wer.dll`, `amsi.dll`, `CRYPTSP.dll`\n\n### Step 7: Linux Persistence Detection\n\n1. **Cron-based persistence** (T1053.003):\n   - Monitor `/var/spool/cron/crontabs/`, `/etc/crontab`, `/etc/cron.d/`, `/etc/cron.{hourly,daily,weekly,monthly}/`\n   - Alert on new cron entries, especially those downloading and executing remote scripts\n   - `auditd` rule: `-w /etc/crontab -p wa -k cron_modification`\n2. **Systemd persistence** (T1543.002):\n   - Monitor `/etc/systemd/system/`, `~/.config/systemd/user/` for new unit files\n   - Alert on `systemctl enable` of unknown services\n   - Check for services with `ExecStart` pointing to unusual paths\n3. **SSH authorized keys** (T1098.004):\n   - Monitor `~/.ssh/authorized_keys` for new entries on all users\n   - Alert on authorized_keys modifications for root and service accounts\n   - `auditd` rule: `-w /root/.ssh/authorized_keys -p wa -k ssh_key_change`\n4. **Shell profile persistence**:\n   - Monitor `/etc/profile`, `/etc/profile.d/`, `~/.bashrc`, `~/.bash_profile`, `~/.zshrc`\n   - Alert on new commands added to shell initialization files\n\n## Automated Baseline Comparison\n\n1. Run Sysinternals Autoruns weekly and compare output against the previous baseline\n2. Use `autorunsc.exe -accepteula -a * -c -h -s -v -vt` to generate a CSV with VirusTotal hash checks\n3. Diff the current output against the last known-good baseline\n4. Investigate any new entries not approved through change management\n\n## Escalation Criteria\n\n- WMI event subscription detected on any system: P1 - Highly indicative of compromise\n- New service or scheduled task running as SYSTEM with unsigned binary: P1\n- Registry run key modifications during off-hours from non-admin accounts: P2\n- SSH authorized key added to a server without corresponding change ticket: P2\n- Multiple persistence mechanisms on the same host: P1 - Indicates advanced adversary establishing redundant access\n\n## References\n\n- MITRE ATT&CK: Persistence (TA0003)\n- Sysinternals Autoruns: https://docs.microsoft.com/en-us/sysinternals/downloads/autoruns\n- SANS: Windows Persistence Mechanism Cheat Sheet\n- MITRE ATT&CK: T1547.001, T1053.005, T1546.003, T1543.003, T1574.001, T1053.003, T1543.002, T1098.004",
  "content_type": "runbook",
  "category": "threat_detection",
  "subcategory": "persistence",
  "tags": ["persistence", "registry-run-keys", "scheduled-tasks", "wmi-subscription", "startup-folder", "services", "dll-hijacking", "cron", "systemd", "autoruns"],
  "severity_filter": ["medium", "high", "critical"],
  "incident_types": ["persistence", "malware", "backdoor", "active_intrusion"],
  "ioc_types": ["hash", "ip", "domain"],
  "mitre_techniques": ["T1547.001", "T1053.005", "T1546.003", "T1543.003", "T1574.001", "T1053.003", "T1543.002", "T1098.004", "T1547.004", "T1547.009"],
  "compliance_frameworks": ["NIST", "SOC2", "ISO27001"],
  "priority": 85
}
