# kb-content-output/articles/threat_detection/privilege-escalation-detection.json
{
  "title": "Privilege Escalation Detection Guide",
  "content": "## Purpose\n\nThis SOP defines the detection strategies for identifying privilege escalation attempts across Windows and Linux environments. Privilege escalation allows adversaries to gain higher-level permissions than initially obtained, enabling access to sensitive resources, lateral movement, and persistence. Detection of privilege escalation is essential for containing compromises before the adversary achieves domain or root-level access.\n\n## Scope\n\nCovers detection of both local and domain-level privilege escalation techniques including token manipulation, UAC bypass, kernel exploits, service account abuse, group policy modification, and sudo/SUID abuse on Linux systems.\n\n## Prerequisites\n\n- Sysmon deployed with a comprehensive configuration (SwiftOnSecurity or Olaf Hartong baseline)\n- Windows Security Event Log audit policies configured per CIS Benchmarks\n- PowerShell Script Block Logging and Module Logging enabled\n- EDR telemetry with process creation, registry modification, and privilege use monitoring\n- Linux auditd configured with privilege-related rules\n- Active Directory audit logs from domain controllers\n\n## Procedure\n\n### Step 1: Detect Token Manipulation (T1134)\n\n1. **Access Token Theft**:\n   - Monitor Sysmon Event ID 10 (Process Access) where the target process is `lsass.exe`\n   - Access mask flags to watch: `0x1010` (PROCESS_QUERY_LIMITED_INFORMATION + PROCESS_VM_READ), `0x1038`, `0x1FFFFF` (PROCESS_ALL_ACCESS)\n   - Alert on non-security tool processes accessing lsass.exe (exclude legitimate AV, EDR agents)\n   - Windows Event ID 4672 (Special Privileges Assigned to New Logon): Monitor for SeDebugPrivilege, SeImpersonatePrivilege, SeAssignPrimaryTokenPrivilege assigned to unexpected accounts\n2. **Token Impersonation**:\n   - Detect tools: `incognito.exe`, Meterpreter `steal_token` and `impersonate_token` commands\n   - Monitor for processes running under one user context making API calls as a different (higher-privileged) user\n   - Sysmon Event ID 1: Look for token manipulation tool filenames or known hashes\n3. **Potato Exploits** (SeImpersonate abuse):\n   - `JuicyPotato`, `RottenPotato`, `SweetPotato`, `PrintSpoofer`, `GodPotato`\n   - These exploit services running as SYSTEM that can be tricked into authenticating to an attacker-controlled named pipe\n   - Detect: Process running as a service account suddenly spawning `cmd.exe` or `powershell.exe` as SYSTEM\n\n### Step 2: Detect UAC Bypass (T1548.002)\n\n1. **Known UAC bypass techniques**:\n   - `fodhelper.exe`: Monitor for registry modification at `HKCU\\Software\\Classes\\ms-settings\\shell\\open\\command` followed by execution of `fodhelper.exe`\n   - `computerdefaults.exe`: Similar registry-based hijack at `HKCU\\Software\\Classes\\ms-settings`\n   - `eventvwr.exe`: Registry modification at `HKCU\\Software\\Classes\\mscfile\\shell\\open\\command`\n   - `sdclt.exe`: Registry hijack at `HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\App Paths\\control.exe`\n2. **Generic UAC bypass detection**:\n   - Sysmon Event ID 13 (Registry Value Set): Monitor for modifications to `HKCU\\Software\\Classes\\*\\shell\\open\\command` followed within seconds by auto-elevated binary execution\n   - Process tree: Auto-elevated Microsoft binary spawning unexpected child processes (cmd.exe, powershell.exe, or unsigned binaries)\n   - Event ID 4688 with Token Elevation Type = 1 (Full Token) from unexpected parent processes\n3. **DLL hijacking for elevation**:\n   - Monitor for DLL writes to directories searched by auto-elevated programs (e.g., `C:\\Windows\\System32\\` by non-installer processes)\n   - Sysmon Event ID 7 (Image Loaded): Unsigned DLLs loaded by elevated processes\n\n### Step 3: Detect Kernel Exploits (T1068)\n\n1. **Known exploit indicators**:\n   - BSOD or system instability followed by new SYSTEM-level processes\n   - Driver loading events (Sysmon Event ID 6) for unsigned or suspicious drivers\n   - Event ID 7045: New kernel-mode service installation\n2. **Bring Your Own Vulnerable Driver (BYOVD)**:\n   - Monitor for loading of known vulnerable drivers (maintain a list from loldrivers.io)\n   - Alert on any new driver load from unusual paths (not `C:\\Windows\\System32\\drivers\\`)\n   - Sysmon Event ID 6 with signature status != Valid\n3. **Exploit mitigation monitoring**:\n   - Check if Credential Guard, HVCI, and Exploit Protection are enabled on targeted systems\n   - Monitor Windows Defender Exploit Guard events for blocked exploit attempts\n\n### Step 4: Detect Service Account Abuse\n\n1. **Service creation for privilege escalation** (T1543.003):\n   - Event ID 7045 (Service Installed): New services running as LocalSystem, NT AUTHORITY\\SYSTEM, or with domain admin credentials\n   - Alert on services with binary paths pointing to temp directories, user profile paths, or writable shares\n   - Services with spaces in the unquoted path (unquoted service path vulnerability)\n2. **Service configuration modification** (T1574.011):\n   - `sc config` commands changing the service binary path or logon account\n   - Registry modifications to `HKLM\\SYSTEM\\CurrentControlSet\\Services\\<name>\\ImagePath`\n3. **Scheduled task privilege escalation** (T1053.005):\n   - Event ID 4698 (Scheduled Task Created): Tasks running as SYSTEM or elevated accounts created by non-admin users\n   - `schtasks /create /sc` with `/ru SYSTEM` from non-admin context\n\n### Step 5: Detect Active Directory Privilege Escalation\n\n1. **Group membership changes**:\n   - Event ID 4728 (Member Added to Security-Enabled Global Group)\n   - Event ID 4732 (Member Added to Security-Enabled Local Group)\n   - Event ID 4756 (Member Added to Security-Enabled Universal Group)\n   - Critical groups: Domain Admins, Enterprise Admins, Schema Admins, Administrators, Account Operators, Backup Operators\n2. **DCSync attack** (T1003.006):\n   - Event ID 4662: Object access on the domain with properties `{1131f6aa-9c07-11d1-f79f-00c04fc2dcd2}` (DS-Replication-Get-Changes) or `{1131f6ad-9c07-11d1-f79f-00c04fc2dcd2}` (DS-Replication-Get-Changes-All)\n   - Source is not a domain controller\n3. **Kerberoasting** (T1558.003):\n   - Event ID 4769: TGS requests with RC4 encryption (0x17) for service accounts\n   - High volume of TGS requests from a single source for different SPNs\n\n### Step 6: Detect Linux Privilege Escalation\n\n1. **Sudo abuse**:\n   - Monitor `/var/log/auth.log` for `sudo` usage by unexpected users\n   - Alert on `sudoers` file modifications (`/etc/sudoers`, `/etc/sudoers.d/`)\n   - Detect `sudo -l` enumeration (attacker checking what they can sudo)\n2. **SUID/SGID exploitation**:\n   - Detect new SUID binaries: `find / -perm -4000 -newer /reference_timestamp`\n   - Alert on SUID binaries in unusual locations (`/tmp`, `/home`, `/var`)\n   - Monitor for known SUID escalation tools: GTFOBins-listed binaries with SUID set\n3. **Kernel exploit indicators**:\n   - Unexpected kernel module loading (`insmod`, `modprobe`)\n   - Process running as unprivileged user suddenly spawning root shells\n   - Core dumps from kernel-level crashes\n\n## Escalation Criteria\n\n- Any confirmed privilege escalation to Domain Admin or Enterprise Admin: Immediate P1 - Treat as full domain compromise\n- DCSync or Kerberoasting detected: P1 - Credential theft in progress\n- UAC bypass on crown jewel systems: P1\n- New SYSTEM-level service or scheduled task from unknown binary: P2 investigation required\n- SUID/sudo abuse on Linux servers: P2 with containment\n- Token manipulation or Potato exploit detected: P1 - Active exploitation\n\n## References\n\n- MITRE ATT&CK: Privilege Escalation (TA0004)\n- MITRE ATT&CK: T1134, T1548.002, T1068, T1543.003, T1053.005, T1003.006, T1558.003\n- LOLDrivers Project: https://www.loldrivers.io/\n- GTFOBins: https://gtfobins.github.io/\n- SANS: Detecting Privilege Escalation with Windows Event Logs",
  "content_type": "sop",
  "category": "threat_detection",
  "subcategory": "privilege_escalation",
  "tags": ["privilege-escalation", "token-manipulation", "uac-bypass", "kernel-exploit", "service-abuse", "dcsync", "kerberoasting", "suid", "sudo"],
  "severity_filter": ["high", "critical"],
  "incident_types": ["privilege_escalation", "credential_theft", "active_intrusion", "apt"],
  "ioc_types": ["hash", "ip"],
  "mitre_techniques": ["T1134", "T1134.001", "T1134.002", "T1548.002", "T1068", "T1543.003", "T1053.005", "T1003.006", "T1558.003", "T1574.011"],
  "compliance_frameworks": ["NIST", "SOC2", "ISO27001", "PCI-DSS"],
  "priority": 88
}
