# kb-content-output/articles/threat_detection/siem-alert-triage.json
{
  "title": "SIEM Alert Triage SOP",
  "content": "## Purpose\n\nThis Standard Operating Procedure defines the structured approach for triaging alerts generated by the Security Information and Event Management (SIEM) platform. Effective alert triage ensures that genuine threats are identified and escalated promptly while reducing analyst fatigue caused by false positives and informational noise.\n\n## Scope\n\nThis SOP applies to all Tier 1 and Tier 2 SOC analysts responsible for monitoring and triaging SIEM alerts across all connected data sources, including endpoint telemetry, network sensors, cloud workloads, identity platforms, and email gateways.\n\n## Prerequisites\n\n- Active access to the SIEM console with analyst-level permissions\n- Familiarity with the organization's asset inventory and network topology\n- Access to threat intelligence feeds (MISP, VirusTotal, AbuseIPDB)\n- Understanding of baseline network behavior and normal business operations\n- Access to the ticketing system for incident creation and tracking\n\n## Procedure\n\n### Step 1: Initial Alert Review\n\n1. Open the SIEM alert queue and sort by severity (Critical > High > Medium > Low).\n2. Review the alert title, description, and the detection rule that fired.\n3. Note the timestamp, source host, destination, user account, and triggering event.\n4. Check if the alert is part of a known campaign or ongoing investigation by searching the case management system.\n\n### Step 2: Contextual Enrichment\n\n1. **Asset Context**: Identify the affected asset in the CMDB. Determine its business criticality (crown jewel, standard workstation, development server, etc.).\n2. **User Context**: Look up the associated user account in the identity provider. Check role, department, recent activity, and whether the account is a service account.\n3. **Threat Intelligence**: Query the IOCs (IPs, domains, hashes) against threat intel platforms. Check for known malicious indicators, reputation scores, and associated threat actor groups.\n4. **Historical Context**: Search the SIEM for prior alerts involving the same source, destination, or user within the past 30 days. Identify patterns or repeated behavior.\n\n### Step 3: Alert Classification\n\nClassify the alert into one of the following categories:\n\n| Classification | Description | Action |\n|---|---|---|\n| **True Positive** | Confirmed malicious or policy-violating activity | Escalate to Tier 2 / create incident |\n| **Benign True Positive** | Legitimate activity correctly detected by the rule (e.g., authorized pen test) | Document and close with annotation |\n| **False Positive** | Non-malicious activity incorrectly flagged | Close alert, submit tuning request |\n| **Insufficient Data** | Cannot determine without additional investigation | Escalate for deeper analysis |\n\n### Step 4: Correlation Analysis\n\n1. Pivot on shared indicators (IP, hostname, user, process hash) to find related alerts within a +/- 2-hour window.\n2. Check for kill-chain progression: Was there a preceding reconnaissance alert? Is there a subsequent lateral movement or data exfiltration alert?\n3. Identify if multiple low-severity alerts correlate to form a higher-severity composite event.\n4. Use MITRE ATT&CK mapping to determine the attack stage:\n   - **Initial Access** (T1566, T1190): Phishing, exploit of public-facing application\n   - **Execution** (T1059, T1204): Script execution, user-initiated malicious file\n   - **Persistence** (T1547, T1053): Registry modification, scheduled task creation\n   - **Lateral Movement** (T1021, T1570): Remote service use, lateral tool transfer\n   - **Exfiltration** (T1041, T1048): Data transfer over C2 channel or alternative protocol\n\n### Step 5: Priority Assignment\n\nAssign priority using the following matrix:\n\n| Asset Criticality | Threat Confidence High | Threat Confidence Medium | Threat Confidence Low |\n|---|---|---|---|\n| **Crown Jewel** | P1 - Critical | P1 - Critical | P2 - High |\n| **Business Critical** | P1 - Critical | P2 - High | P3 - Medium |\n| **Standard** | P2 - High | P3 - Medium | P4 - Low |\n| **Development/Test** | P3 - Medium | P4 - Low | P4 - Low |\n\n### Step 6: Documentation and Handoff\n\n1. Create or update the incident ticket with all enrichment findings, classification rationale, and assigned priority.\n2. Include raw log excerpts, screenshots of relevant SIEM queries, and IOC listings.\n3. If escalating, provide a brief analyst narrative summarizing the alert, context gathered, and recommended next steps.\n4. Tag the ticket with relevant MITRE ATT&CK technique IDs and affected asset identifiers.\n\n## Escalation Criteria\n\n- **Immediate Escalation (Tier 2 / Incident Commander)**:\n  - Any P1 alert confirmed as True Positive\n  - Evidence of active data exfiltration or ransomware encryption\n  - Compromise of domain admin or service account credentials\n  - Alerts involving crown jewel assets with medium or higher confidence\n\n- **Standard Escalation (Tier 2)**:\n  - P2 alerts requiring deeper forensic analysis\n  - Correlated alert chains suggesting multi-stage attack\n  - Alerts classified as Insufficient Data after 30 minutes of triage\n\n- **Tuning Request (Detection Engineering)**:\n  - Repeated false positives from the same rule (3+ in 7 days)\n  - Rules generating excessive volume without actionable context\n\n## Key Performance Indicators\n\n- Mean Time to Triage (MTTT): Target < 15 minutes per alert\n- False Positive Rate: Target < 30% per detection rule\n- Escalation Accuracy: Target > 90% of escalated alerts confirmed as true positives\n\n## References\n\n- NIST SP 800-61 Rev. 2 - Computer Security Incident Handling Guide\n- MITRE ATT&CK Framework: https://attack.mitre.org\n- SANS Incident Handler's Handbook\n- Internal Asset Classification Policy\n- SOC Shift Handoff Procedure",
  "content_type": "sop",
  "category": "threat_detection",
  "subcategory": "alert_triage",
  "tags": ["siem", "alert-triage", "soc", "classification", "correlation", "false-positive", "priority-assignment"],
  "severity_filter": ["low", "medium", "high", "critical"],
  "incident_types": ["generic_alert", "multi_stage_attack", "unauthorized_access", "malware"],
  "ioc_types": ["ip", "domain", "hash", "url", "email"],
  "mitre_techniques": ["T1566", "T1190", "T1059", "T1204", "T1547", "T1053", "T1021", "T1570", "T1041", "T1048"],
  "compliance_frameworks": ["NIST", "SOC2", "ISO27001"],
  "priority": 95
}
