# kb-content-output/articles/identity_access/privileged-session-monitoring.json
{
  "title": "Privileged Session Monitoring and Recording",
  "content": "# Privileged Session Monitoring and Recording\n\n## Purpose\n\nThis procedure defines the requirements and operational processes for monitoring and recording privileged sessions -- any interactive session where a user exercises elevated privileges on critical infrastructure. Privileged session monitoring (PSM) provides real-time visibility into administrative actions, creates an audit trail for compliance and forensics, and enables detection of insider threats and compromised admin accounts.\n\n## Scope\n\nAll privileged sessions including: domain administrator sessions on Active Directory domain controllers, root/sudo sessions on Linux servers, database administrator sessions (DBA), cloud management console sessions with admin roles, network device administrative sessions (routers, switches, firewalls), hypervisor/VMware vCenter admin sessions, and application admin sessions for critical business applications.\n\n## Privileged Session Categories\n\n| Category | Systems | Risk Level | Monitoring Requirement |\n|----------|---------|------------|------------------------|\n| Tier 0 | Domain Controllers, CA servers, Azure AD Connect, IdP | Critical | Full session recording + real-time monitoring |\n| Tier 1 | Member servers, databases, cloud admin consoles | High | Full session recording |\n| Tier 2 | Workstations with local admin, developer environments | Medium | Audit logging, selective recording |\n| Network | Firewalls, routers, switches, load balancers | High | Command logging + session recording |\n| Cloud | AWS/Azure/GCP management consoles with admin roles | High | API logging + console session recording |\n\n## Session Monitoring Architecture\n\n### Components\n\n1. **Privileged Access Management (PAM) system**: CyberArk, BeyondTrust, Delinea, HashiCorp Boundary\n2. **Session proxy/gateway**: Routes all privileged sessions through a monitored jump host\n3. **Session recorder**: Captures full keystroke/screen activity\n4. **Real-time analytics**: Detects anomalous commands or behaviors during active sessions\n5. **Audit log aggregation**: Centralizes session metadata in SIEM\n\n### Session Flow\n\n```\nAdmin User -> MFA Challenge -> PAM Portal -> Credential Checkout -> Session Proxy -> Target System\n                                                                        |\n                                                                   [Recording]\n                                                                   [Keystroke Log]\n                                                                   [Real-time Alerts]\n```\n\n## Session Events to Monitor\n\n### Linux Privileged Sessions\n\n```bash\n# Configure auditd for privileged command monitoring\n# /etc/audit/rules.d/privileged.rules\n\n# Monitor all commands run as root\n-a always,exit -F arch=b64 -F euid=0 -S execve -k privileged-commands\n-a always,exit -F arch=b32 -F euid=0 -S execve -k privileged-commands\n\n# Monitor sudo usage\n-w /etc/sudoers -p wa -k sudoers-change\n-w /etc/sudoers.d/ -p wa -k sudoers-change\n\n# Monitor su usage\n-w /bin/su -p x -k su-execution\n\n# Monitor SSH key changes\n-w /root/.ssh/ -p wa -k ssh-root-keys\n-w /home/ -p wa -k ssh-user-keys\n\n# Monitor passwd/shadow changes\n-w /etc/passwd -p wa -k identity-file\n-w /etc/shadow -p wa -k identity-file\n-w /etc/group -p wa -k identity-file\n\n# Monitor cron changes\n-w /etc/crontab -p wa -k cron-change\n-w /var/spool/cron/ -p wa -k cron-change\n```\n\n```bash\n# Query auditd for privileged commands\nausearch -k privileged-commands --start today -i | head -100\n\n# Find all sudo commands executed today\nausearch -k privileged-commands -i --start today | grep -E \"type=EXECVE\" | head -50\n\n# Monitor sudo sessions in real time\ntail -f /var/log/auth.log | grep -E \"sudo|su\\[\"\n```\n\n### Windows Privileged Sessions\n\n```kql\n// Detect interactive logons to Domain Controllers\nSecurityEvent\n| where EventID == 4624\n| where Computer in (domain_controllers)\n| where LogonType in (2, 10, 11)  // Interactive, RemoteInteractive, CachedInteractive\n| where Account !endswith \"$\"  // Exclude machine accounts\n| project TimeGenerated, Account, Computer, LogonType, IpAddress, WorkstationName\n| order by TimeGenerated desc\n\n// Detect PowerShell execution on sensitive servers\nEvent\n| where EventLog == \"Microsoft-Windows-PowerShell/Operational\"\n| where EventID == 4104  // Script block logging\n| where Computer in (sensitive_servers)\n| extend ScriptBlock = tostring(EventData)\n| where ScriptBlock has_any (\"Invoke-Mimikatz\", \"Get-ADReplAccount\", \"Set-DomainObject\",\n        \"Add-DomainGroupMember\", \"Invoke-Command\", \"New-PSSession\", \"Enter-PSSession\")\n| project TimeGenerated, Computer, ScriptBlock\n\n// Detect RDP sessions to sensitive servers\nSecurityEvent\n| where EventID == 4624\n| where LogonType == 10  // RemoteInteractive (RDP)\n| where Computer in (sensitive_servers)\n| project TimeGenerated, Account, Computer, IpAddress\n| order by TimeGenerated desc\n```\n\n### Database Privileged Sessions\n\n```sql\n-- PostgreSQL - Enable and monitor privileged session logging\n-- postgresql.conf settings:\n-- log_connections = on\n-- log_disconnections = on\n-- log_statement = 'ddl'  -- Log all DDL statements\n-- log_line_prefix = '%m [%p] %q%u@%d '\n\n-- Query PostgreSQL logs for DDL operations\nSELECT log_time, user_name, database_name, command_tag, message\nFROM postgres_log\nWHERE command_tag IN ('CREATE TABLE', 'DROP TABLE', 'ALTER TABLE',\n                      'CREATE ROLE', 'DROP ROLE', 'ALTER ROLE',\n                      'GRANT', 'REVOKE', 'COPY', 'TRUNCATE')\nORDER BY log_time DESC;\n\n-- MySQL/MariaDB - Monitor privileged operations\nSELECT event_time, user_host, command_type, argument\nFROM mysql.general_log\nWHERE command_type IN ('Query', 'Execute')\n  AND argument REGEXP 'GRANT|REVOKE|CREATE USER|DROP USER|ALTER USER|FLUSH PRIVILEGES'\nORDER BY event_time DESC;\n```\n\n### Network Device Sessions\n\n```bash\n# Configure syslog for network device command logging\n\n# Cisco IOS - Enable command logging\n# archive\n#  log config\n#   logging enable\n#   logging size 1000\n#   notify syslog contenttype plaintext\n#   hidekeys\n\n# Parse network device logs for privileged commands\ngrep -E \"(configure terminal|enable|reload|write|copy|crypto|access-list|route-map)\" /var/log/network-devices.log | tail -50\n\n# Alert on configuration changes\ngrep -E \"(SYS-5-CONFIG_I|CONFIG_CHANGE|CFGLOG)\" /var/log/network-devices.log\n```\n\n### Cloud Console Session Monitoring\n\n```kql\n// AWS - Monitor console sign-ins with admin access\n// CloudTrail: ConsoleLogin events for admin users\nAWSCloudTrail\n| where EventName == \"ConsoleLogin\"\n| where ResponseElements has \"Success\"\n| extend UserArn = tostring(UserIdentity.arn)\n| where UserArn has \"admin\" or UserArn has \"AdministratorAccess\"\n| project TimeGenerated, UserArn, SourceIPAddress, UserAgent\n| order by TimeGenerated desc\n\n// Azure - Monitor portal sessions with admin roles\nSigninLogs\n| where AppDisplayName == \"Azure Portal\"\n| where ResultType == 0\n| join kind=inner (\n    AuditLogs\n    | where OperationName == \"Add member to role\"\n    | extend RoleName = tostring(TargetResources[0].modifiedProperties[1].newValue)\n    | where RoleName has_any (\"Global Administrator\", \"Privileged Role Administrator\",\n            \"Security Administrator\", \"Exchange Administrator\")\n    | extend AdminUser = tostring(TargetResources[0].userPrincipalName)\n    | project AdminUser, RoleName\n  ) on $left.UserPrincipalName == $right.AdminUser\n| project TimeGenerated, UserPrincipalName, RoleName, IPAddress, Location\n```\n\n## Anomalous Session Detection\n\n### Behavioral Analytics Rules\n\n| Detection Rule | Threshold | Action |\n|---------------|-----------|--------|\n| Admin login outside business hours | Before 6 AM or after 10 PM local time | Alert SOC |\n| Admin login from new IP/location | First-time IP for that admin | Alert SOC |\n| Bulk data query by DBA | > 100K rows in single query | Alert SOC |\n| Multiple failed sudo attempts | > 5 failures in 10 minutes | Alert SOC + temporary lockout |\n| Session duration exceeds normal | > 2x typical session length | Alert SOC |\n| Privilege escalation chain | su/sudo to root, then network connection | Alert SOC immediately |\n| Configuration file modification | Changes to /etc/, registry, GPO | Alert SOC |\n| New scheduled task/cron during session | Task persistence creation | Alert SOC immediately |\n| Data exfiltration indicators | Large file copy, compression, SCP/SFTP | Alert SOC immediately |\n\n### Real-Time Command Alerting\n\nConfigure real-time alerts for high-risk commands during privileged sessions:\n\n```bash\n# High-risk Linux commands to alert on in real-time:\n# - useradd, usermod, userdel (account manipulation)\n# - chmod 777, chmod +s (permission changes)\n# - iptables -F, ufw disable (firewall modification)\n# - history -c, unset HISTFILE (anti-forensics)\n# - dd if=/dev/sda (disk imaging)\n# - curl/wget to unknown domains (downloading tools)\n# - nc -l, ncat (reverse shells)\n# - python -c, perl -e, ruby -e (scripting engines)\n# - base64 -d, openssl enc (decoding/decrypting payloads)\n# - mount -o remount,rw (filesystem modification)\n# - systemctl disable, service stop (security service tampering)\n```\n\n## Session Recording Requirements\n\n| Requirement | Tier 0 | Tier 1 | Tier 2 |\n|------------|--------|--------|--------|\n| Full video/screen recording | Required | Required | Optional |\n| Keystroke logging | Required | Required | Audit log only |\n| Command capture | Required | Required | Required |\n| File transfer logging | Required | Required | Optional |\n| Clipboard monitoring | Required | Optional | Not required |\n| Real-time monitoring | Required | Alerting only | Not required |\n| Recording retention | 1 year | 1 year | 90 days |\n| Tamper protection | Required | Required | Required |\n\n## Compliance Mapping\n\n| Framework | Requirement | PSM Control |\n|-----------|-------------|-------------|\n| PCI-DSS 4.0 | 10.2.2 - Log admin actions | Session recording + audit logs |\n| SOC2 | CC6.1 - Logical access controls | Session monitoring + alerting |\n| HIPAA | 164.312(b) - Audit controls | Session logging for ePHI systems |\n| NIST 800-53 | AC-2(12) - Account monitoring | Real-time session monitoring |\n| ISO 27001 | A.9.4 - System access control | Privileged access management |\n| SOX | Section 404 - Internal controls | Session recording for financial systems |\n\n## Escalation Criteria\n\n- Unauthorized privileged session on Tier 0 system: escalate to IR Lead immediately\n- Anti-forensics commands detected (history clearing, log deletion): escalate to IR Lead\n- Privileged session from unknown IP or location: escalate to SOC Lead\n- Bulk data access by DBA outside approved maintenance window: escalate to SOC Lead and Data Owner\n- Session recording gap or tamper detected: escalate to PAM engineering team\n- Multiple admin accounts used from same source in short timeframe: escalate to IR Lead (possible credential dump)\n\n## MITRE ATT&CK Mapping\n\n- **T1078** -- Valid Accounts\n- **T1078.002** -- Valid Accounts: Domain Accounts\n- **T1021** -- Remote Services\n- **T1021.001** -- Remote Services: Remote Desktop Protocol\n- **T1021.004** -- Remote Services: SSH\n- **T1059** -- Command and Scripting Interpreter\n- **T1070.003** -- Indicator Removal: Clear Command History\n- **T1098** -- Account Manipulation\n\n## References\n\n- NIST SP 800-53 Rev 5: AC-2(12), AU-14 Session Audit\n- CIS Controls v8: Control 5 - Account Management, Control 6 - Access Control Management\n- PCI-DSS v4.0: Requirement 10 - Log and Monitor All Access\n- SANS: Privileged Access Management Best Practices\n- CyberArk: Privileged Session Management Documentation\n- BeyondTrust: Privileged Remote Access Guide",
  "content_type": "procedure",
  "category": "identity_access",
  "subcategory": "privileged_access",
  "tags": ["privileged-access", "session-monitoring", "pam", "session-recording", "rdp", "ssh", "sudo", "admin-session", "insider-threat", "auditd", "database-admin"],
  "severity_filter": ["medium", "high", "critical"],
  "incident_types": ["insider_threat", "unauthorized_access", "privilege_abuse", "account_compromise", "data_exfiltration"],
  "ioc_types": ["ip"],
  "mitre_techniques": ["T1078", "T1078.002", "T1021", "T1021.001", "T1021.004", "T1059", "T1070.003", "T1098"],
  "compliance_frameworks": ["NIST", "SOC2", "ISO27001", "PCI-DSS", "HIPAA", "SOX"],
  "priority": 73
}