# kb-content-output/articles/cloud_security/container-security.json
{
  "title": "Container and Kubernetes Security Monitoring",
  "content": "# Container and Kubernetes Security Monitoring\n\n## Purpose\n\nThis procedure defines the security monitoring approach for containerized workloads and Kubernetes clusters. It covers pod security policies, container image scanning, runtime threat detection, and network policy enforcement across managed Kubernetes services (EKS, AKS, GKE) and self-hosted clusters.\n\n## Scope\n\nAll Kubernetes clusters and container runtimes under SOC monitoring, including production, staging, and development environments. Production clusters warrant continuous monitoring; non-production clusters require at minimum daily review.\n\n## Container Image Security\n\n### Image Scanning Pipeline\n\nAll container images must be scanned before deployment and continuously in registries:\n\n**Pre-deployment scanning (CI/CD pipeline):**\n\n1. **Static vulnerability scanning**: Scan images for known CVEs using Trivy, Grype, or Snyk Container\n```bash\n# Trivy scan with severity filter\ntrivy image --severity HIGH,CRITICAL --exit-code 1 <registry>/<image>:<tag>\n\n# Grype scan\ngrype <registry>/<image>:<tag> --fail-on high\n```\n\n2. **Configuration scanning**: Check Dockerfiles for misconfigurations\n```bash\n# Check for common Dockerfile security issues\n# - Running as root (no USER instruction)\n# - Using latest tag instead of pinned version\n# - Including secrets in build arguments\n# - Unnecessary packages installed\ntrivy config --severity HIGH,CRITICAL Dockerfile\n```\n\n3. **Secret scanning**: Detect embedded secrets in image layers\n```bash\n# Scan image layers for secrets\ntrivy image --scanners secret <registry>/<image>:<tag>\n```\n\n**Registry scanning (continuous):**\n- Enable automatic scanning in container registries (ECR, ACR, Artifact Registry, Harbor)\n- Configure alerts for newly discovered CVEs in deployed images\n- Set policies to block deployment of images with critical/high CVEs older than 30 days\n\n### Image Policy Enforcement\n\nImplement admission controllers to enforce image security policies:\n\n```yaml\n# OPA/Gatekeeper ConstraintTemplate - Require images from trusted registries\napiVersion: templates.gatekeeper.sh/v1\nkind: ConstraintTemplate\nmetadata:\n  name: k8sallowedrepos\nspec:\n  crd:\n    spec:\n      names:\n        kind: K8sAllowedRepos\n      validation:\n        openAPIV3Schema:\n          type: object\n          properties:\n            repos:\n              type: array\n              items:\n                type: string\n  targets:\n    - target: admission.k8s.gatekeeper.sh\n      rego: |\n        package k8sallowedrepos\n        violation[{\"msg\": msg}] {\n          container := input.review.object.spec.containers[_]\n          satisfied := [good | repo = input.parameters.repos[_]; good = startswith(container.image, repo)]\n          not any(satisfied)\n          msg := sprintf(\"container '%v' has an invalid image repo '%v'\", [container.name, container.image])\n        }\n```\n\n## Kubernetes Security Monitoring\n\n### Audit Log Monitoring\n\nKubernetes API server audit logs are the primary source for detecting malicious activity. Monitor for:\n\n**Critical Events (Alert Immediately):**\n\n| Event | Description | Risk |\n|-------|-------------|------|\n| Privileged pod creation | `spec.securityContext.privileged: true` | Container escape |\n| Host path mount | `spec.volumes[].hostPath` on sensitive paths | Node compromise |\n| Host network/PID namespace | `spec.hostNetwork: true` or `spec.hostPID: true` | Network sniffing, process manipulation |\n| ClusterRole binding to cluster-admin | RBAC granting full cluster access | Privilege escalation |\n| ServiceAccount token mount in non-system namespace | Default SA token available in pod | Lateral movement |\n| Secret access from unexpected pods | GET/LIST on secrets resources | Credential theft |\n| exec into production pods | `pods/exec` subresource | Interactive access, potential data theft |\n\n**Suspicious Events (Review within 1 hour):**\n- DaemonSet creation (runs on every node)\n- Namespace creation (potential staging for attack)\n- ConfigMap/Secret creation with encoded data\n- CronJob creation (persistence mechanism)\n- Role/RoleBinding modifications\n- Node labeling changes (can affect scheduling)\n\n### Falco Runtime Detection Rules\n\nDeploy Falco or equivalent runtime security tools for real-time container behavior monitoring:\n\n```yaml\n# Key Falco rules for SOC monitoring\n- rule: Terminal shell in container\n  desc: Detect interactive shell sessions in containers\n  condition: >\n    spawned_process and container and\n    proc.name in (bash, sh, zsh, dash) and\n    proc.tty != 0\n  output: \"Shell spawned in container (user=%user.name container=%container.name shell=%proc.name)\"\n  priority: WARNING\n\n- rule: Contact K8S API from container\n  desc: Detect container communicating with Kubernetes API\n  condition: >\n    outbound and container and\n    fd.sip.name = \"kubernetes.default.svc.cluster.local\"\n  output: \"Container accessing K8s API (container=%container.name image=%container.image.repository)\"\n  priority: NOTICE\n\n- rule: Crypto mining binary detected\n  desc: Detect known crypto mining process names\n  condition: >\n    spawned_process and container and\n    proc.name in (xmrig, minerd, minergate, cpuminer)\n  output: \"Crypto miner detected (container=%container.name process=%proc.name)\"\n  priority: CRITICAL\n\n- rule: Sensitive file read\n  desc: Detect reading of sensitive files in containers\n  condition: >\n    open_read and container and\n    fd.name in (/etc/shadow, /etc/passwd, /root/.ssh/authorized_keys, /root/.bash_history)\n  output: \"Sensitive file read in container (file=%fd.name container=%container.name)\"\n  priority: WARNING\n```\n\n### Network Policy Monitoring\n\nKubernetes NetworkPolicies provide microsegmentation. Monitor for:\n\n1. **Missing network policies**: Namespaces without any NetworkPolicy (all traffic allowed by default)\n2. **Overly permissive policies**: Policies allowing all ingress/egress\n3. **Policy violations**: Traffic blocked by NetworkPolicies (indicates misconfiguration or attack)\n\n```bash\n# Find namespaces without NetworkPolicies\nfor ns in $(kubectl get ns -o jsonpath='{.items[*].metadata.name}'); do\n  count=$(kubectl get networkpolicies -n $ns --no-headers 2>/dev/null | wc -l)\n  if [ \"$count\" -eq 0 ]; then\n    echo \"WARNING: No NetworkPolicy in namespace: $ns\"\n  fi\ndone\n```\n\n**Recommended baseline NetworkPolicy:**\n```yaml\n# Default deny all ingress and egress\napiVersion: networking.k8s.io/v1\nkind: NetworkPolicy\nmetadata:\n  name: default-deny-all\nspec:\n  podSelector: {}\n  policyTypes:\n    - Ingress\n    - Egress\n```\n\n## Managed Kubernetes Service Monitoring\n\n### EKS (AWS)\n- Enable EKS audit logs to CloudWatch Logs\n- Enable GuardDuty EKS Protection for runtime monitoring\n- Monitor EKS add-on and configuration changes via CloudTrail\n\n### AKS (Azure)\n- Enable diagnostic settings for kube-audit and kube-audit-admin logs\n- Enable Microsoft Defender for Containers\n- Monitor AKS configuration changes via Azure Activity Log\n\n### GKE (Google Cloud)\n- Enable GKE audit logging (admin activity is always on; enable data access)\n- Enable GKE Threat Detection in Security Command Center\n- Monitor GKE configuration changes via Cloud Audit Logs\n\n## Escalation Criteria\n\nEscalate to Incident Response when:\n- Evidence of container escape (host namespace access, privileged container with host mounts)\n- Crypto mining detected in production clusters\n- Unauthorized cluster-admin access granted\n- Secrets exfiltration from Kubernetes\n- Malicious container image deployed (backdoor, web shell, reverse shell)\n- Supply chain compromise (base image or dependency tampered)\n\n## MITRE ATT&CK Mapping (Containers)\n\n- **T1610** — Deploy Container\n- **T1611** — Escape to Host\n- **T1613** — Container and Resource Discovery\n- **T1552.007** — Container API\n- **T1053.007** — Scheduled Task/Job: Container Orchestration Job\n\n## References\n\n- [NSA/CISA Kubernetes Hardening Guide](https://media.defense.gov/2022/Aug/29/2003066362/-1/-1/0/CTR_KUBERNETES_HARDENING_GUIDANCE_1.2_20220829.PDF)\n- [Falco Rules Reference](https://falco.org/docs/reference/rules/default-rules/)\n- [CIS Kubernetes Benchmark](https://www.cisecurity.org/benchmark/kubernetes)",
  "content_type": "procedure",
  "category": "cloud_security",
  "subcategory": "containers",
  "tags": ["kubernetes", "containers", "docker", "eks", "aks", "gke", "falco", "runtime-security", "image-scanning", "network-policy"],
  "severity_filter": ["medium", "high", "critical"],
  "incident_types": ["container_compromise", "cryptomining", "unauthorized_access", "supply_chain_attack"],
  "ioc_types": ["ip", "domain", "hash"],
  "mitre_techniques": ["T1610", "T1611", "T1613", "T1552.007", "T1053.007"],
  "compliance_frameworks": ["NIST", "SOC2", "ISO27001", "PCI-DSS"],
  "priority": 50
}