# kb-content-output/articles/endpoint_security/endpoint-detection-gaps.json
{
  "title": "Identifying Endpoint Detection Gaps",
  "content": "## Purpose\n\nThis procedure provides a systematic methodology for identifying gaps in endpoint detection coverage. Detection gaps are blind spots where adversary techniques can execute without generating alerts. Regularly assessing and remediating these gaps is essential for maintaining an effective security posture. This guide covers MITRE ATT&CK coverage mapping, telemetry gap analysis, and sensor health monitoring.\n\n## Scope\n\nApplies to SOC engineering teams, detection engineers, and security architects responsible for maintaining and improving endpoint detection capabilities across the EDR, SIEM, and endpoint logging infrastructure.\n\n## Prerequisites\n\n- MITRE ATT&CK Navigator access (https://mitre-attack.github.io/attack-navigator/)\n- Inventory of all detection rules (SIEM rules, EDR behavioral rules, custom analytics)\n- Inventory of deployed endpoint agents and their configurations\n- Endpoint telemetry data source documentation\n- Access to EDR management console for agent health monitoring\n- Recent threat intelligence reports relevant to the organization's industry\n\n## Procedure\n\n### Step 1: Map Current Detection Coverage to MITRE ATT&CK\n\n1. **Inventory all detection rules**:\n   - Export all SIEM detection rules with their MITRE ATT&CK technique mappings\n   - Export all EDR behavioral detection rules and their technique coverage\n   - Include custom analytics, UEBA rules, and threat intelligence matching rules\n   - Include threat hunting queries that are run regularly (even if not automated)\n\n2. **Create the ATT&CK Navigator heatmap**:\n   - For each technique, assess coverage level:\n     - **Full Detection** (score: 100, green): Automated detection rule with proven true positive rate > 70%\n     - **Partial Detection** (score: 50, yellow): Detection exists but with known limitations (high FP rate, limited data source, or only detects specific sub-techniques)\n     - **Hunt Only** (score: 25, light yellow): No automated detection, but included in regular threat hunting queries\n     - **No Detection** (score: 0, red): No automated or manual detection capability\n   - Save the navigator layer as the baseline for tracking improvement over time\n\n3. **Calculate coverage metrics**:\n   - Total techniques in scope: Filter ATT&CK techniques applicable to your environment (Windows, Linux, cloud, etc.)\n   - Coverage rate: (Techniques with Full or Partial Detection / Total In-Scope Techniques) x 100\n   - Target: > 60% full/partial coverage for endpoint-relevant techniques\n\n### Step 2: Identify Critical Gaps\n\nNot all gaps have equal risk. Prioritize based on threat relevance.\n\n1. **Threat-informed prioritization**:\n   - Review the top 10 threat actors targeting your industry (use MITRE ATT&CK Groups or CTI reports)\n   - Map their known TTPs to your coverage heatmap\n   - Any technique used by a relevant threat actor that you cannot detect is a critical gap\n\n2. **Kill-chain gap analysis**:\n   - For each phase of the kill chain, verify you have at least one detection:\n\n   | Kill Chain Phase | ATT&CK Tactic | Minimum Detection Requirement |\n   |---|---|---|\n   | Delivery | Initial Access | Email attachment, malicious link, exploit detection |\n   | Exploitation | Execution | Process creation, script execution monitoring |\n   | Installation | Persistence | Autoruns, service creation, scheduled task monitoring |\n   | C2 | Command and Control | Network anomaly, DNS monitoring, proxy analysis |\n   | Actions on Objectives | Exfiltration/Impact | DLP, volume anomaly, ransomware behavior |\n\n   - A gap at any single phase means an attack can progress undetected through that stage\n\n3. **Common blind spots to check**:\n   - **Defense Evasion (TA0005)**: Often the weakest coverage area. Check for:\n     - Process injection detection (T1055): Requires Sysmon Event 8 or EDR behavioral rules\n     - AMSI bypass detection (T1562.001): Requires PowerShell Script Block Logging\n     - Timestomping detection (T1070.006): Requires file integrity monitoring\n     - Log clearing detection (T1070.001): Requires Event ID 1102 forwarding\n   - **Credential Access (TA0006)**:\n     - LSASS access monitoring (T1003.001): Requires Sysmon Event 10\n     - Kerberoasting detection (T1558.003): Requires Kerberos ticket request analysis\n     - DCSync detection (T1003.006): Requires directory service access auditing\n   - **Collection (TA0009)**:\n     - Screen capture (T1113): Rarely detected by endpoint tools\n     - Clipboard data (T1115): Requires specialized monitoring\n     - Email collection (T1114): Requires email platform auditing\n\n### Step 3: Telemetry Gap Analysis\n\nDetection rules are only as good as the data they analyze. Verify that required telemetry is being collected.\n\n1. **Map data sources per technique**:\n\n   | ATT&CK Technique | Required Data Source | Status |\n   |---|---|---|\n   | T1059.001 PowerShell | Script Block Logging (4104) | Verify enabled |\n   | T1055 Process Injection | Sysmon Event 8 (CreateRemoteThread) | Verify Sysmon config |\n   | T1003.001 LSASS Dump | Sysmon Event 10 (ProcessAccess) | Verify LSASS is monitored |\n   | T1021.002 SMB/Admin Share | Event 5145 (Network Share Access) | Verify audit policy |\n   | T1543.003 Service Creation | Event 7045 (System Log) | Verify forwarding |\n   | T1053.005 Scheduled Task | Event 4698 (Task Created) | Verify audit policy |\n\n2. **Verify data is reaching the SIEM**:\n   - For each critical data source, run a simple query to confirm recent data exists\n   - Check for gaps: Look for periods where expected data is missing\n   - Verify field parsing: Ensure the SIEM is correctly extracting fields needed for detection rules\n\n3. **Common telemetry gaps**:\n   - **Sysmon not deployed or misconfigured**: Sysmon provides 26 event types; default configuration may not enable all needed types\n   - **Command line logging not enabled**: Event ID 4688 without command line is nearly useless for investigation. Verify GPO setting.\n   - **PowerShell logging not enabled**: Script Block Logging (4104) and Module Logging (4103) must be explicitly enabled\n   - **DNS query logging not configured**: Sysmon Event 22 or Windows DNS debug logging required\n   - **Object access auditing not enabled**: Events 4663, 4656, 5145 require explicit SACL configuration\n\n### Step 4: Sensor Health Monitoring\n\n1. **EDR agent health**:\n   - Monitor the percentage of endpoints with healthy, active EDR agents\n   - Target: > 98% of managed endpoints reporting within the last 24 hours\n   - Alert on endpoints that stop reporting (agent crashed, disabled, or uninstalled)\n   - Track agent version: Endpoints on outdated agent versions may lack latest detection capabilities\n\n2. **SIEM data source health**:\n   - Create dashboards showing event volume per data source over time\n   - Set minimum event volume thresholds: If Windows Security Event Log volume drops below 50% of the 7-day average, investigate\n   - Monitor log source latency: Events arriving more than 5 minutes late may cause time-sensitive correlation rules to miss events\n\n3. **Configuration drift detection**:\n   - Periodically audit Sysmon configurations across the fleet to verify consistency\n   - Check Group Policy application: `gpresult /r` on sample endpoints to verify audit policies are applied\n   - Verify EDR policy assignment: Ensure all endpoints are assigned the correct detection policy\n\n### Step 5: Gap Remediation Planning\n\n1. **Prioritize remediation by risk**:\n   - **Priority 1**: Gaps in techniques used by threat actors actively targeting your industry\n   - **Priority 2**: Gaps in the kill chain that leave an entire attack phase undetected\n   - **Priority 3**: Telemetry gaps that affect multiple detection rules\n   - **Priority 4**: Coverage improvements for techniques with no known active exploitation\n\n2. **Remediation approaches**:\n   - **Enable telemetry**: Deploy Sysmon, enable audit policies, configure logging (often the quickest fix)\n   - **Write new detection rules**: Create Sigma rules or SIEM-specific detections for uncovered techniques\n   - **Tune existing rules**: Improve existing rules that have partial coverage\n   - **Deploy new tools**: If no existing tool can provide the needed telemetry (e.g., network detection for encrypted C2)\n   - **Add to hunt program**: For techniques that are difficult to detect reliably with automation, add them to the regular threat hunting program\n\n3. **Track remediation progress**:\n   - Update the ATT&CK Navigator heatmap monthly\n   - Track coverage rate trend over time (should be increasing)\n   - Set quarterly goals: \"Improve coverage of TA0005 Defense Evasion from 35% to 55%\"\n\n### Step 6: Validation with Purple Team Exercises\n\n1. **Simulate uncovered techniques**:\n   - Use Atomic Red Team (https://github.com/redcanaryco/atomic-red-team) to safely simulate ATT&CK techniques\n   - For each gap identified, run the corresponding atomic test\n   - After deploying a new detection rule, validate it detects the atomic test\n\n2. **Scheduled validation cadence**:\n   - Monthly: Test 5-10 highest priority techniques\n   - Quarterly: Comprehensive validation of all critical technique detections\n   - After major changes: Validate detection coverage after SIEM upgrades, EDR agent updates, or infrastructure changes\n\n3. **Document validation results**:\n   - For each technique tested: Pass/Fail, detection latency, alert fidelity\n   - Failed validations become remediation tickets\n\n## Escalation Criteria\n\n- EDR agent coverage drops below 95%: Escalate to endpoint management team for immediate remediation\n- Critical data source stops forwarding to SIEM: Escalate to SIEM engineering as P1 operational issue\n- No detection coverage for techniques used in active campaign targeting your industry: Escalate to security leadership with risk assessment and remediation timeline\n- Detection validation reveals a previously unknown gap in a critical kill chain phase: Escalate to detection engineering for priority remediation\n\n## References\n\n- MITRE ATT&CK Navigator: https://mitre-attack.github.io/attack-navigator/\n- MITRE ATT&CK Data Sources: https://attack.mitre.org/datasources/\n- Atomic Red Team: https://github.com/redcanaryco/atomic-red-team\n- Sigma Rules: https://github.com/SigmaHQ/sigma\n- NIST SP 800-53 Rev. 5: CA-7 (Continuous Monitoring), SI-4 (System Monitoring)\n- DeTT&CT Framework: https://github.com/rabobank-cdc/DeTTECT",
  "content_type": "procedure",
  "category": "endpoint_security",
  "subcategory": "detection_engineering",
  "tags": ["detection-gaps", "mitre-attack", "coverage-mapping", "telemetry", "sensor-health", "purple-team", "atomic-red-team", "sysmon", "detection-engineering"],
  "severity_filter": ["low", "medium", "high", "critical"],
  "incident_types": [],
  "ioc_types": [],
  "mitre_techniques": ["T1055", "T1562.001", "T1070.001", "T1070.006", "T1003.001", "T1003.006", "T1558.003", "T1059.001"],
  "compliance_frameworks": ["NIST", "SOC2", "ISO27001"],
  "priority": 70
}
