# kb-content-output/articles/endpoint_security/linux-endpoint-investigation.json
{
  "title": "Linux Endpoint Investigation Guide",
  "content": "## Purpose\n\nThis runbook provides step-by-step procedures for investigating potentially compromised Linux endpoints. Linux systems frequently serve as web servers, database hosts, container platforms, and cloud infrastructure, making them high-value targets. This guide covers authentication analysis, process investigation, file integrity verification, persistence review, and rootkit detection.\n\n## Scope\n\nApplicable to Linux distributions commonly used in enterprise environments including Ubuntu, Debian, RHEL, CentOS, Amazon Linux, and Alpine. Covers investigation via SSH, EDR remote shell, or direct console access.\n\n## Prerequisites\n\n- SSH access to the target system or EDR remote shell capability\n- Root or sudo access for full investigation capabilities\n- Forensic toolkit transferred to the system via trusted media (not downloaded from the internet on the compromised host)\n- External storage or network share for evidence export\n- Understanding of the system's role (web server, database, container host, etc.)\n\n## Runbook Steps\n\n### Step 1: Authentication Log Analysis\n\nAuthentication logs are the primary source for identifying unauthorized access.\n\n**Log locations**:\n- Debian/Ubuntu: `/var/log/auth.log`\n- RHEL/CentOS: `/var/log/secure`\n- Systemd journal: `journalctl -u sshd --since \"7 days ago\"`\n\n**Key patterns to search for**:\n\n1. **Successful SSH logins**:\n   ```\n   grep \"Accepted\" /var/log/auth.log\n   ```\n   Review:\n   - Logins from unexpected source IPs\n   - Logins at unusual times (off-hours for the system's timezone)\n   - Key-based vs. password-based authentication (key-based with unknown key = compromised key)\n   - Logins to root account directly (should be disabled)\n\n2. **Failed SSH attempts**:\n   ```\n   grep \"Failed password\" /var/log/auth.log | awk '{print $(NF-3)}' | sort | uniq -c | sort -rn | head -20\n   ```\n   This extracts source IPs with the most failed attempts. High counts indicate brute force.\n\n3. **Invalid user attempts**:\n   ```\n   grep \"Invalid user\" /var/log/auth.log\n   ```\n   Attackers often try common usernames (admin, test, oracle, postgres, deploy).\n\n4. **Sudo usage audit**:\n   ```\n   grep \"sudo:\" /var/log/auth.log\n   ```\n   Review:\n   - Users running sudo commands they should not need\n   - Sudo failures (\"user NOT in sudoers\" = potential enumeration)\n   - Sudo to root shell: `sudo -i`, `sudo bash`, `sudo su -`\n\n5. **Account changes**:\n   ```\n   grep -E \"useradd|usermod|userdel|groupadd|passwd\" /var/log/auth.log\n   ```\n   Any account creation or modification during the suspected compromise period is suspicious.\n\n### Step 2: Process Investigation\n\n1. **List all running processes with full details**:\n   ```\n   ps auxf --cols 300 > /evidence/process_list.txt\n   ```\n   Review for:\n   - Processes running from `/tmp`, `/dev/shm`, `/var/tmp`, `/run` (world-writable directories)\n   - Processes with suspicious names mimicking system processes (e.g., `[kworker/0:1]` with brackets in the name used to mimic kernel threads)\n   - Processes consuming unusual CPU or memory\n   - Processes running as root that should not be (web applications, user tools)\n\n2. **Check for deleted executables** (fileless or anti-forensic malware):\n   ```\n   ls -la /proc/*/exe 2>/dev/null | grep deleted\n   ```\n   A process whose executable has been deleted from disk is highly suspicious. The binary is still in memory but was removed to hinder analysis.\n\n3. **Recover deleted executables**:\n   ```\n   cp /proc/<PID>/exe /evidence/recovered_binary\n   sha256sum /evidence/recovered_binary\n   ```\n\n4. **Check for processes with network connections**:\n   ```\n   ss -tupna | grep ESTAB\n   ```\n   Correlate PIDs with the process list. Investigate:\n   - Outbound connections to unknown external IPs\n   - Listening services on unexpected ports\n   - Connections to known C2 infrastructure\n\n5. **Check for hidden processes**:\n   - Compare `ps aux` output with `/proc` directory listing:\n   ```\n   diff <(ps -eo pid --no-headers | sort -n) <(ls /proc | grep -E '^[0-9]+$' | sort -n)\n   ```\n   Discrepancies may indicate a rootkit hiding processes.\n\n### Step 3: File System Investigation\n\n1. **Find recently modified files**:\n   ```\n   find / -type f -mtime -7 -not -path '/proc/*' -not -path '/sys/*' -not -path '/run/*' 2>/dev/null | head -500\n   ```\n   Focus on files modified in system directories (`/usr/bin`, `/usr/sbin`, `/lib`, `/etc`).\n\n2. **Find SUID/SGID binaries** (potential privilege escalation):\n   ```\n   find / -type f \\( -perm -4000 -o -perm -2000 \\) -exec ls -la {} \\; 2>/dev/null > /evidence/suid_sgid.txt\n   ```\n   Compare against a known-good baseline. New SUID binaries are a strong indicator of compromise.\n\n3. **Check for web shells** (if web server is running):\n   ```\n   find /var/www /srv/www /opt -name '*.php' -mtime -30 2>/dev/null\n   find /var/www /srv/www /opt -name '*.jsp' -mtime -30 2>/dev/null\n   ```\n   Review new or recently modified script files in web-accessible directories.\n   Common web shell indicators in PHP:\n   - `eval(`, `base64_decode(`, `system(`, `exec(`, `passthru(`, `shell_exec(`\n   - `$_GET['cmd']`, `$_POST['cmd']`, `$_REQUEST['cmd']`\n\n4. **Check temporary directories**:\n   ```\n   ls -la /tmp/ /dev/shm/ /var/tmp/\n   ```\n   Attackers frequently stage tools and output in world-writable directories.\n\n5. **Check for unauthorized SSH keys**:\n   ```\n   find / -name authorized_keys -exec echo \"=== {} ===\" \\; -exec cat {} \\; 2>/dev/null\n   ```\n   Any unauthorized keys provide persistent backdoor access.\n\n### Step 4: Persistence Review\n\n1. **Cron jobs** (all locations):\n   ```\n   echo \"=== System crontab ===\"\n   cat /etc/crontab\n   echo \"=== Cron.d ===\"\n   ls -la /etc/cron.d/\n   cat /etc/cron.d/*\n   echo \"=== User crontabs ===\"\n   for user in $(cut -f1 -d: /etc/passwd); do\n     crontab_output=$(crontab -u $user -l 2>/dev/null)\n     if [ -n \"$crontab_output\" ]; then\n       echo \"=== $user ===\"\n       echo \"$crontab_output\"\n     fi\n   done\n   ```\n   Look for: Download and execute patterns, reverse shells, base64-encoded commands.\n\n2. **Systemd services**:\n   ```\n   systemctl list-unit-files --type=service --state=enabled\n   find /etc/systemd/system/ -name \"*.service\" -newer /etc/os-release\n   ```\n   Review new or modified service unit files. Check `ExecStart` for suspicious commands.\n\n3. **Init scripts and rc.local**:\n   ```\n   cat /etc/rc.local\n   ls -la /etc/init.d/\n   ```\n\n4. **Shell profile backdoors**:\n   ```\n   cat /etc/profile\n   cat /etc/bash.bashrc\n   ls -la /etc/profile.d/\n   for dir in /home/* /root; do\n     for file in .bashrc .bash_profile .profile .bash_login .zshrc; do\n       if [ -f \"$dir/$file\" ]; then\n         echo \"=== $dir/$file ===\"\n         cat \"$dir/$file\"\n       fi\n     done\n   done\n   ```\n\n5. **LD_PRELOAD hijacking**:\n   ```\n   cat /etc/ld.so.preload 2>/dev/null\n   env | grep LD_PRELOAD\n   ```\n   Any library listed in `ld.so.preload` is loaded into every dynamically-linked process.\n\n### Step 5: Rootkit Detection\n\n1. **Run chkrootkit**:\n   ```\n   chkrootkit -q\n   ```\n   Review output for infections. Note: Must run the tool from trusted media, not from the potentially compromised system's package.\n\n2. **Run rkhunter**:\n   ```\n   rkhunter --check --skip-keypress --report-warnings-only\n   ```\n\n3. **Manual rootkit indicators**:\n   - Hidden files in system directories: `find /usr /bin /sbin -name \".*\" -type f`\n   - Modified system binaries: Compare SHA256 of `/usr/bin/ps`, `/usr/bin/ls`, `/usr/bin/netstat` against known-good packages\n     ```\n     debsums -c  # Debian/Ubuntu: Check installed packages for modified files\n     rpm -Va     # RHEL/CentOS: Verify all installed packages\n     ```\n   - Kernel module rootkits: `lsmod` and compare against expected modules. Check for hidden modules:\n     ```\n     diff <(cat /proc/modules | awk '{print $1}' | sort) <(lsmod | awk '{print $1}' | sort)\n     ```\n\n### Step 6: Network Configuration Review\n\n1. **Check network interfaces for promiscuous mode** (sniffing):\n   ```\n   ip link | grep PROMISC\n   ```\n\n2. **Review iptables/nftables rules**:\n   ```\n   iptables -L -n -v\n   nft list ruleset\n   ```\n   Look for rules that redirect traffic, allow unexpected inbound access, or masquerade traffic.\n\n3. **DNS configuration**:\n   ```\n   cat /etc/resolv.conf\n   ```\n   Verify DNS servers have not been changed to attacker-controlled resolvers.\n\n4. **hosts file**:\n   ```\n   cat /etc/hosts\n   ```\n   Check for entries redirecting security update domains or C2 domains.\n\n## Escalation Criteria\n\n- Rootkit detected: P1 - System cannot be trusted; reimage from known-good media\n- Unauthorized SSH keys on root or service accounts: P1 - Persistent backdoor access\n- Deleted executable still running (fileless malware): P1 - Active compromise\n- Modified system binaries (ls, ps, netstat): P1 - Rootkit or advanced compromise\n- Web shell found in production web server: P1 - Active web server compromise\n- Unauthorized user accounts created: P2 - Investigate scope and timeline\n\n## References\n\n- NIST SP 800-86: Guide to Integrating Forensic Techniques into Incident Response\n- SANS Linux Forensics Cheat Sheet\n- chkrootkit: http://www.chkrootkit.org/\n- rkhunter: https://rkhunter.sourceforge.net/\n- MITRE ATT&CK: T1059.004 (Unix Shell), T1053.003 (Cron), T1543.002 (Systemd Service), T1098.004 (SSH Authorized Keys), T1014 (Rootkit)",
  "content_type": "runbook",
  "category": "endpoint_security",
  "subcategory": "linux_investigation",
  "tags": ["linux", "investigation", "auth-logs", "cron", "ssh", "rootkit", "process-analysis", "file-integrity", "web-shell"],
  "severity_filter": ["medium", "high", "critical"],
  "incident_types": ["malware", "active_intrusion", "unauthorized_access", "web_shell", "rootkit"],
  "ioc_types": ["ip", "hash", "domain"],
  "mitre_techniques": ["T1059.004", "T1053.003", "T1543.002", "T1098.004", "T1014", "T1505.003", "T1078.003", "T1070.002", "T1222.002"],
  "compliance_frameworks": ["NIST", "SOC2", "ISO27001", "PCI-DSS"],
  "priority": 80
}
