# kb-content-output/articles/endpoint_security/macos-endpoint-investigation.json
{
  "title": "macOS Endpoint Investigation Guide",
  "content": "# macOS Endpoint Investigation Guide\n\n## Overview\n\nThis guide provides SOC analysts with macOS-specific forensic investigation techniques, artifact locations, and tool usage for incident response. macOS has unique security mechanisms (Gatekeeper, XProtect, TCC, SIP) and artifact locations that differ significantly from Windows and Linux. Understanding these differences is essential for effective investigation of macOS endpoints.\n\nmacOS is based on Darwin (BSD-derived) and uses APFS (Apple File System) by default. Investigations must account for APFS snapshots, sealed system volumes (SSV), and Apple's privacy protections (TCC) that may restrict forensic tool access.\n\n## Key macOS Security Mechanisms\n\n| Mechanism | Description | Investigation Relevance |\n|-----------|-------------|------------------------|\n| Gatekeeper | Verifies code signing and notarization for downloaded applications | Bypass indicates sophisticated attacker or user override |\n| XProtect | Built-in signature-based malware detection | Check XProtect logs for known malware detection |\n| MRT (Malware Removal Tool) | Removes known malware families | Check if MRT has run and what it removed |\n| SIP (System Integrity Protection) | Prevents modification of protected system files | Disabling SIP is a strong indicator of compromise |\n| TCC (Transparency, Consent, and Control) | Controls application access to sensitive data | Unauthorized TCC modifications indicate malware |\n| Secure Boot | Verifies boot chain integrity | Boot policy changes may indicate rootkit |\n| FileVault | Full-disk encryption | Must be unlocked for live forensics |\n\n## Critical Artifact Locations\n\n### System Logs\n\n| Log | Path | Content |\n|-----|------|--------|\n| Unified Log | `/var/db/diagnostics/` and `/var/db/uuidtext/` | All system, application, and security events |\n| Install Log | `/var/log/install.log` | Software installations and updates |\n| System Log (legacy) | `/var/log/system.log` | System events (being replaced by Unified Log) |\n| Audit Logs | `/var/audit/` | BSM audit trail (if enabled) |\n| Crash Reports | `~/Library/Logs/DiagnosticReports/` | Application crashes, may contain exploit evidence |\n| XProtect Remediator | `/var/log/xprotect_remediator.log` | XProtect scan and remediation actions |\n\n### Persistence Locations\n\n| Persistence Type | Path | Command to Check |\n|-----------------|------|------------------|\n| LaunchDaemons (system) | `/Library/LaunchDaemons/` | `ls -la /Library/LaunchDaemons/` |\n| LaunchAgents (system) | `/Library/LaunchAgents/` | `ls -la /Library/LaunchAgents/` |\n| LaunchAgents (user) | `~/Library/LaunchAgents/` | `ls -la ~/Library/LaunchAgents/` |\n| Login Items | `~/Library/Application Support/com.apple.backgroundtaskmanagementagent/backgrounditems.btm` | `sfltool dumpbtm` (macOS 13+) |\n| Cron Jobs | `/var/at/tabs/` and `/etc/crontab` | `crontab -l` |\n| Periodic Scripts | `/etc/periodic/daily/`, `/etc/periodic/weekly/`, `/etc/periodic/monthly/` | `ls -la /etc/periodic/*/` |\n| Login/Logout Hooks | `/var/root/Library/Preferences/com.apple.loginwindow.plist` | `defaults read com.apple.loginwindow` |\n| Authorization Plugins | `/Library/Security/SecurityAgentPlugins/` | `ls -la /Library/Security/SecurityAgentPlugins/` |\n| Kernel Extensions | `/Library/Extensions/` | `kextstat` |\n| System Extensions | `/Library/SystemExtensions/` | `systemextensionsctl list` |\n| Profiles | `/var/db/ConfigurationProfiles/` | `profiles list` |\n\n### User Activity Artifacts\n\n| Artifact | Path | Content |\n|----------|------|--------|\n| Bash History | `~/.bash_history` | Command history (bash) |\n| Zsh History | `~/.zsh_history` | Command history (zsh, default since Catalina) |\n| Downloads | `~/Downloads/` | Downloaded files (check quarantine attributes) |\n| Quarantine Events DB | `~/Library/Preferences/com.apple.LaunchServices.QuarantineEventsV2` | Record of all downloaded files with source URL |\n| Recent Items | `~/Library/Application Support/com.apple.sharedfilelist/` | Recently opened files, applications, servers |\n| Spotlight Metadata | `/.Spotlight-V100/` | File metadata index |\n| Finder Preferences | `~/Library/Preferences/com.apple.finder.plist` | Recent folders, sidebar items |\n| Safari History | `~/Library/Safari/History.db` | Web browsing history |\n| Chrome History | `~/Library/Application Support/Google/Chrome/Default/History` | Chrome browsing history |\n| Keychain | `~/Library/Keychains/` | Stored credentials and certificates |\n| TCC Database | `~/Library/Application Support/com.apple.TCC/TCC.db` | Application permission grants |\n\n### Network Artifacts\n\n| Artifact | Path | Content |\n|----------|------|--------|\n| Known WiFi Networks | `/Library/Preferences/com.apple.wifi.known-networks.plist` | Previously connected WiFi networks |\n| DNS Cache | Live query: `sudo dscacheutil -cachedump -entries` | Current DNS cache |\n| Firewall Config | `/usr/libexec/ApplicationFirewall/socketfilterfw` | Application firewall rules |\n| Network Interfaces | Live query: `ifconfig -a` | Interface configuration |\n| ARP Cache | Live query: `arp -a` | Current ARP table |\n| Active Connections | Live query: `netstat -an` or `lsof -i` | Current network connections |\n\n## Investigation Procedures\n\n### Step 1: Unified Log Analysis\n\nThe Unified Log is the primary log source on modern macOS. Use the `log` command to query:\n\n```bash\n# Show all logs for the last hour\nlog show --last 1h\n\n# Filter by process name\nlog show --predicate 'process == \"loginwindow\"' --last 24h\n\n# Search for specific event types\nlog show --predicate 'eventMessage CONTAINS \"authentication\"' --last 24h\n\n# Filter by subsystem (security events)\nlog show --predicate 'subsystem == \"com.apple.securityd\"' --last 24h\n\n# Look for sudo usage\nlog show --predicate 'process == \"sudo\"' --last 7d\n\n# Look for SSH connections\nlog show --predicate 'process == \"sshd\"' --last 7d\n\n# Look for Gatekeeper events\nlog show --predicate 'subsystem == \"com.apple.syspolicy\"' --last 7d\n\n# Look for TCC access decisions\nlog show --predicate 'subsystem == \"com.apple.TCC\"' --last 7d\n\n# Export logs to a file for offline analysis\nlog collect --output ~/Desktop/system_logs.logarchive --last 7d\n```\n\n### Step 2: Persistence Investigation\n\n```bash\n# List all LaunchDaemons and LaunchAgents\necho \"=== System LaunchDaemons ===\"\nls -la /Library/LaunchDaemons/\necho \"=== System LaunchAgents ===\"\nls -la /Library/LaunchAgents/\necho \"=== User LaunchAgents ===\"\nls -la ~/Library/LaunchAgents/\n\n# Check for recently modified plist files (last 7 days)\nfind /Library/LaunchDaemons /Library/LaunchAgents ~/Library/LaunchAgents -name \"*.plist\" -mtime -7 2>/dev/null\n\n# Verify code signing on launch items\nfor plist in /Library/LaunchDaemons/*.plist; do\n  program=$(/usr/libexec/PlistBuddy -c \"Print :Program\" \"$plist\" 2>/dev/null || /usr/libexec/PlistBuddy -c \"Print :ProgramArguments:0\" \"$plist\" 2>/dev/null)\n  if [ -n \"$program\" ] && [ -f \"$program\" ]; then\n    echo \"--- $plist ---\"\n    codesign -dvv \"$program\" 2>&1 | grep -E \"Authority|TeamIdentifier|Identifier\"\n  fi\ndone\n\n# Check cron jobs\ncrontab -l 2>/dev/null\nfor user in $(dscl . -list /Users | grep -v '^_'); do\n  echo \"=== Cron for $user ===\"\n  crontab -u \"$user\" -l 2>/dev/null\ndone\n\n# Check login items (macOS 13+)\nsfltool dumpbtm 2>/dev/null\n\n# Check for kernel extensions (rare on modern macOS)\nkextstat | grep -v com.apple\n\n# Check for system extensions\nsystemextensionsctl list\n\n# Check for configuration profiles\nprofiles list -verbose\n```\n\n### Step 3: Process and Application Analysis\n\n```bash\n# List all running processes with full paths\nps -eo pid,ppid,user,command\n\n# Check for unsigned or ad-hoc signed processes\nfor pid in $(ps -eo pid=); do\n  path=$(ps -p $pid -o comm= 2>/dev/null)\n  if [ -n \"$path\" ] && [ -f \"$path\" ]; then\n    sig=$(codesign -dvv \"$path\" 2>&1)\n    if echo \"$sig\" | grep -q \"code object is not signed\"; then\n      echo \"UNSIGNED: PID=$pid PATH=$path\"\n    elif echo \"$sig\" | grep -q \"adhoc\"; then\n      echo \"AD-HOC SIGNED: PID=$pid PATH=$path\"\n    fi\n  fi\ndone\n\n# List all open files for a specific process\nlsof -p <PID>\n\n# List all network connections by process\nlsof -i -P -n\n\n# Check for recently installed applications\nls -lt /Applications/ | head -20\nmdfind 'kMDItemContentType == \"com.apple.application-bundle\"' -onlyin /Applications/ | while read app; do\n  echo \"$(stat -f '%Sm' -t '%Y-%m-%d %H:%M' \"$app\") $app\"\ndone | sort -r | head -20\n\n# Check application quarantine status\nxattr -l /Applications/*.app 2>/dev/null | grep quarantine\n\n# Examine running launch items\nlaunchctl list | grep -v com.apple\n```\n\n### Step 4: User Activity Investigation\n\n```bash\n# Check shell history\ncat ~/.zsh_history 2>/dev/null | tail -100\ncat ~/.bash_history 2>/dev/null | tail -100\n\n# Check Quarantine Events database (downloaded file history)\nsqlite3 ~/Library/Preferences/com.apple.LaunchServices.QuarantineEventsV2 \\\n  \"SELECT datetime(LSQuarantineTimeStamp + 978307200, 'unixepoch') as date, \\\n  LSQuarantineAgentName, LSQuarantineDataURLString, LSQuarantineOriginURLString \\\n  FROM LSQuarantineEvent ORDER BY LSQuarantineTimeStamp DESC LIMIT 50;\"\n\n# Check TCC database for permission grants\nsqlite3 ~/Library/Application\\ Support/com.apple.TCC/TCC.db \\\n  \"SELECT client, service, auth_value, auth_reason, last_modified FROM access ORDER BY last_modified DESC LIMIT 50;\"\n\n# Check for recently modified files in user directories\nfind ~/ -type f -mtime -1 -not -path '*/Library/Caches/*' -not -path '*/.Trash/*' 2>/dev/null | head -50\n\n# Check SSH known_hosts and authorized_keys\ncat ~/.ssh/known_hosts 2>/dev/null\ncat ~/.ssh/authorized_keys 2>/dev/null\n\n# Check for recently used sudo\nlog show --predicate 'process == \"sudo\"' --last 24h --style compact\n```\n\n### Step 5: Network Investigation\n\n```bash\n# Current network connections\nnetstat -an | grep ESTABLISHED\nlsof -i -P -n | grep ESTABLISHED\n\n# DNS cache\nsudo dscacheutil -cachedump -entries Host 2>/dev/null\n\n# Check firewall status\nsudo /usr/libexec/ApplicationFirewall/socketfilterfw --getglobalstate\nsudo /usr/libexec/ApplicationFirewall/socketfilterfw --listapps\n\n# Check for VPN configurations\nscutil --nc list\n\n# Check hosts file for modifications\ncat /etc/hosts\n\n# Check DNS resolver configuration\nscutil --dns\n\n# Known WiFi networks\ndefaults read /Library/Preferences/com.apple.wifi.known-networks 2>/dev/null\n```\n\n### Step 6: Malware and Threat Indicators\n\n```bash\n# Check if SIP is enabled\ncsrutil status\n\n# Check XProtect definitions version\nsystem_profiler SPInstallHistoryDataType | grep -A 2 \"XProtect\"\n\n# Check for known malware persistence paths\nfind /tmp /var/tmp ~/Library/Application\\ Support -name \"*.dylib\" -o -name \"*.so\" 2>/dev/null\n\n# Look for suspicious hidden files\nfind / -name \".*\" -type f -not -path \"/System/*\" -not -path \"/Library/.*\" -newer /etc/passwd 2>/dev/null | head -50\n\n# Check for DYLD environment variable abuse\nlaunchctl print system | grep -i dyld\nenv | grep DYLD\n\n# Scan for unsigned dylibs loaded in running processes\nfor pid in $(ps -eo pid=); do\n  vmmap $pid 2>/dev/null | grep -i dylib | grep -v /usr/lib | grep -v /System\ndone\n\n# Check for recently modified binaries in common paths\nfind /usr/local/bin /usr/local/sbin -type f -mtime -7 2>/dev/null\n```\n\n## macOS-Specific Attack Techniques\n\n| MITRE Technique | macOS Implementation | Detection |\n|----------------|---------------------|----------|\n| T1543.004 | Launch Daemon/Agent persistence | Monitor plist creation in LaunchDaemons/LaunchAgents directories |\n| T1547.015 | Login Items persistence | Monitor BTM database and sfltool output |\n| T1059.004 | Bash/Zsh script execution | Shell history, Unified Log process events |\n| T1059.002 | AppleScript/JXA execution | Monitor osascript process execution |\n| T1222.002 | File permission modification | Monitor chmod/chown on sensitive files |\n| T1564.001 | Hidden files and directories | Scan for dotfiles in unexpected locations |\n| T1553.001 | Gatekeeper bypass | Unified Log syspolicy events, xattr removal |\n| T1548.004 | Elevated execution via AuthorizationExecuteWithPrivileges | Monitor Security.framework usage |\n| T1056.002 | GUI input capture | TCC database monitoring for input monitoring permissions |\n| T1555.001 | Keychain credential access | Monitor security command-line tool usage, keychain access events |\n| T1574.004 | Dylib hijacking | Monitor DYLD_INSERT_LIBRARIES, unusual dylib loads |\n| T1562.001 | SIP bypass/disable | csrutil status monitoring, boot policy changes |\n\n## Forensic Collection Tools\n\n| Tool | Purpose | Usage |\n|------|---------|-------|\n| log (built-in) | Unified Log collection and analysis | `log collect --output /path/to/output.logarchive` |\n| system_profiler | Hardware and software inventory | `system_profiler SPSoftwareDataType SPNetworkDataType` |\n| ditto | Forensic file copy preserving metadata | `ditto --rsrc /source /destination` |\n| tmutil | Time Machine snapshot management | `tmutil listlocalsnapshots /` |\n| pmset | Power management and wake event history | `pmset -g log` |\n| sysdiagnose | Comprehensive system diagnostic collection | `sudo sysdiagnose -f /path/to/output/` |\n| fs_usage | Real-time file system activity monitoring | `sudo fs_usage -f filesys` |\n| dtrace | Dynamic tracing (requires SIP partially disabled) | System call tracing |\n| osquery | Cross-platform endpoint visibility | SQL-based system state queries |\n| Crescendo | Real-time macOS event monitoring | Swift-based real-time event tool |\n| KnockKnock | Persistence item scanner | Scans all known persistence locations |\n\n## Escalation Criteria\n\n- SIP disabled on an endpoint: P1 -- Possible rootkit or advanced compromise\n- Unsigned LaunchDaemon installed: P1 -- Likely malware persistence\n- TCC database modified outside System Preferences: P1 -- Privacy bypass\n- XProtect definitions outdated (>30 days): P2 -- Reduced malware protection\n- Unknown kernel extension loaded: P1 -- Kernel-level compromise\n- Gatekeeper bypass detected: P2 -- Unsigned code execution\n\n## References\n\n- Apple Platform Security Guide: https://support.apple.com/guide/security/\n- macOS Security Compliance Project: https://github.com/usnistgov/macos_security\n- Objective-See macOS Security Tools: https://objective-see.org/tools.html\n- MITRE ATT&CK for macOS: https://attack.mitre.org/matrices/enterprise/macos/\n- The Art of Mac Malware: https://taomm.org/",
  "content_type": "procedure",
  "category": "endpoint_security",
  "subcategory": "macos_forensics",
  "tags": ["macos", "endpoint-forensics", "unified-log", "persistence", "launchdaemon", "launchagent", "gatekeeper", "sip", "tcc", "apple-security"],
  "severity_filter": ["low", "medium", "high", "critical"],
  "incident_types": ["malware", "unauthorized_access", "persistence", "credential_theft"],
  "ioc_types": ["hash", "domain", "ip", "file_path"],
  "mitre_techniques": ["T1543.004", "T1547.015", "T1059.004", "T1059.002", "T1553.001", "T1548.004", "T1555.001", "T1574.004", "T1562.001"],
  "compliance_frameworks": ["NIST", "SOC2", "ISO27001"],
  "priority": 72
}