# kb-content-output/articles/endpoint_security/persistence-investigation.json
{
  "title": "Endpoint Persistence Investigation",
  "content": "## Purpose\n\nThis procedure guides SOC analysts through a systematic investigation of endpoint persistence mechanisms during incident response. When a compromise is confirmed or suspected, identifying all persistence mechanisms is essential for complete remediation. Missing even one backdoor allows the adversary to regain access after initial cleanup.\n\n## Scope\n\nCovers hands-on investigation of persistence on Windows and Linux endpoints using both automated tools and manual inspection techniques. Focuses on Autoruns analysis, service enumeration, scheduled task review, registry analysis, and boot-level persistence.\n\n## Prerequisites\n\n- EDR with remote shell or live response capability\n- Sysinternals Autoruns (autorunsc.exe for command-line, Autoruns.exe for GUI)\n- KAPE or similar triage collection toolkit\n- Baseline of known-good Autoruns output for comparison\n- Access to file reputation services (VirusTotal)\n- Administrative access to the target endpoint (local or remote)\n\n## Procedure\n\n### Step 1: Automated Persistence Scan with Autoruns\n\nAutoruns is the most comprehensive tool for enumerating Windows persistence locations.\n\n1. **Run Autoruns from external media** (do not install on the target):\n   ```\n   autorunsc.exe -accepteula -a * -c -h -s -v -vt > autoruns_output.csv\n   ```\n   Flags:\n   - `-a *`: Check all autorun locations\n   - `-c`: CSV output for analysis\n   - `-h`: Show file hashes\n   - `-s`: Verify digital signatures\n   - `-v`: Verify VirusTotal status\n   - `-vt`: Submit unknown hashes to VirusTotal\n\n2. **Analyze the output**:\n   - Sort by VirusTotal detection count: Any entry with detections > 0 is suspicious\n   - Filter for unsigned entries: `Verified == \"Not Verified\"`\n   - Filter for entries in unusual locations: Not in `C:\\Windows\\` or `C:\\Program Files\\`\n   - Compare against the baseline: Any entry not in the known-good list requires investigation\n\n3. **Focus areas in Autoruns**:\n   - **Logon tab**: Programs that run at user logon (Run keys, Startup folder)\n   - **Scheduled Tasks tab**: All scheduled tasks\n   - **Services tab**: Windows services and their binary paths\n   - **Drivers tab**: Kernel-mode drivers\n   - **Winlogon tab**: Winlogon notification DLLs\n   - **Image Hijacks tab**: IFEO (Image File Execution Options) debugger entries\n   - **Boot Execute tab**: Programs that run during boot (before Windows starts)\n   - **KnownDLLs tab**: DLL override entries\n   - **WMI tab**: WMI event subscriptions\n\n### Step 2: Service Enumeration and Analysis\n\n1. **List all services with their binary paths**:\n   ```\n   # PowerShell\n   Get-WmiObject win32_service | Select-Object Name, DisplayName, State, StartMode, PathName | Export-Csv services.csv\n   \n   # Command line\n   sc query type= service state= all > all_services.txt\n   wmic service get name,displayname,pathname,startmode,state /format:csv > services.csv\n   ```\n\n2. **Identify suspicious services**:\n   - Services with binary paths in unusual locations (`%TEMP%`, `C:\\Users\\`, `C:\\ProgramData\\`)\n   - Services with unquoted paths containing spaces (unquoted service path vulnerability)\n   - Services pointing to `cmd.exe /c` or `powershell.exe` as the binary\n   - Services with short random names (e.g., `svc8a3f`, `xagt`)\n   - Services mimicking legitimate service names with subtle changes (e.g., `svchost` vs `svch0st`)\n   - Services with `Auto` start mode that are not recognized\n\n3. **Verify service binaries**:\n   - Check digital signature: `sigcheck.exe -e -v <binary_path>`\n   - Check hash against VirusTotal\n   - Check file creation date: Recently created binaries on a system that has not had software installed are suspicious\n\n### Step 3: Scheduled Task Deep Dive\n\n1. **Export all scheduled tasks**:\n   ```\n   schtasks /query /fo CSV /v > all_tasks.csv\n   ```\n   Or for more detail, examine the XML task definitions:\n   ```\n   Get-ScheduledTask | ForEach-Object { Export-ScheduledTask -TaskName $_.TaskName -TaskPath $_.TaskPath } > all_tasks_xml.txt\n   ```\n\n2. **Review task XML for suspicious elements**:\n   ```xml\n   <!-- Suspicious: Task running as SYSTEM executing a script -->\n   <Principals>\n     <Principal><UserId>S-1-5-18</UserId></Principal>\n   </Principals>\n   <Exec>\n     <Command>powershell.exe</Command>\n     <Arguments>-ep bypass -file C:\\Users\\Public\\update.ps1</Arguments>\n   </Exec>\n   ```\n\n3. **Check hidden tasks**:\n   - Tasks can be hidden from `schtasks /query` by setting the `Hidden` attribute or by using specific security descriptors\n   - Check the file system directly: `dir C:\\Windows\\System32\\Tasks\\ /s /b`\n   - Check the registry: `HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Schedule\\TaskCache\\Tasks\\`\n\n### Step 4: Registry Persistence Analysis\n\n1. **Check all Run key locations** (see Persistence Mechanism Detection Guide for full list)\n2. **Image File Execution Options (IFEO)** (T1546.012):\n   ```\n   reg query \"HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Image File Execution Options\" /s\n   ```\n   - Look for `Debugger` values: `calc.exe` with `Debugger = C:\\malware\\backdoor.exe` means the backdoor runs whenever calc is launched\n   - `GlobalFlag` and `SilentProcessExit` entries can trigger code on process termination\n\n3. **AppInit_DLLs** (T1546.010):\n   ```\n   reg query \"HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Windows\" /v AppInit_DLLs\n   reg query \"HKLM\\SOFTWARE\\Wow6432Node\\Microsoft\\Windows NT\\CurrentVersion\\Windows\" /v AppInit_DLLs\n   ```\n   - Any DLL listed here is loaded into every user-mode process that loads user32.dll\n\n4. **Shell extensions and COM hijacking** (T1546.015):\n   - `HKCU\\Software\\Classes\\CLSID\\` entries that override system COM objects\n   - Unusual InprocServer32 values pointing to non-standard DLLs\n\n5. **Winlogon entries**:\n   ```\n   reg query \"HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon\" /v Shell\n   reg query \"HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon\" /v Userinit\n   ```\n   - `Shell` should be `explorer.exe`; any additional entries are suspicious\n   - `Userinit` should be `C:\\Windows\\system32\\userinit.exe,`; additional entries indicate persistence\n\n### Step 5: Boot and Pre-OS Persistence\n\n1. **Boot Execute** (runs before Windows fully loads):\n   ```\n   reg query \"HKLM\\SYSTEM\\CurrentControlSet\\Control\\Session Manager\" /v BootExecute\n   ```\n   - Normal value: `autocheck autochk *`\n   - Any additional entries are highly suspicious\n\n2. **Bootkit indicators**:\n   - Check the MBR/VBR for modifications (requires specialized tools)\n   - Check UEFI firmware integrity if the platform supports it\n   - Modern platforms with Secure Boot should detect most bootkits\n\n3. **Driver-level persistence** (T1547.006):\n   - Review `HKLM\\SYSTEM\\CurrentControlSet\\Services\\` for entries with `Type = 1` (kernel driver)\n   - Check drivers in `C:\\Windows\\System32\\drivers\\` for unsigned or recently created files\n   - Validate against the known-good driver list\n\n### Step 6: Linux Persistence Investigation\n\n1. **Cron jobs** (all locations):\n   ```\n   cat /etc/crontab\n   ls -la /etc/cron.d/\n   ls -la /etc/cron.{hourly,daily,weekly,monthly}/\n   for user in $(cut -f1 -d: /etc/passwd); do echo \"=== $user ===\"; crontab -u $user -l 2>/dev/null; done\n   ```\n\n2. **Systemd services and timers**:\n   ```\n   systemctl list-unit-files --type=service --state=enabled\n   systemctl list-timers --all\n   find /etc/systemd/system/ /run/systemd/system/ ~/.config/systemd/user/ -name \"*.service\" -newer /reference_date\n   ```\n\n3. **SSH authorized keys** (check ALL users):\n   ```\n   find / -name authorized_keys -exec ls -la {} \\; -exec cat {} \\;\n   ```\n\n4. **Shell profile backdoors**:\n   ```\n   cat /etc/profile\n   cat /etc/bash.bashrc\n   find /etc/profile.d/ -name \"*.sh\" -newer /reference_date\n   for user_home in /home/* /root; do cat $user_home/.bashrc $user_home/.bash_profile $user_home/.profile 2>/dev/null; done\n   ```\n\n5. **Init scripts** (SysVinit):\n   ```\n   ls -la /etc/init.d/\n   ls -la /etc/rc.local\n   ```\n\n6. **LD_PRELOAD and library injection**:\n   ```\n   cat /etc/ld.so.preload\n   echo $LD_PRELOAD\n   ldconfig -p | grep -i unusual_library\n   ```\n\n## Escalation Criteria\n\n- Multiple persistence mechanisms on the same host: P1 - Indicates sophisticated adversary with backup access\n- Kernel-level persistence (driver, bootkit): P1 - Requires specialized remediation (likely reimage)\n- WMI event subscription found: P1 - Rarely used by legitimate software\n- Persistence mechanism present on domain controller: P1 - Domain-level compromise\n- Persistence found on hosts not previously identified as compromised: Expand incident scope, P1\n\n## References\n\n- MITRE ATT&CK: Persistence (TA0003)\n- Sysinternals Autoruns: https://docs.microsoft.com/en-us/sysinternals/downloads/autoruns\n- SANS DFIR: Windows Forensic Analysis Poster\n- Persistence techniques reference: https://persistence-info.github.io/\n- MITRE ATT&CK: T1547, T1053, T1546, T1543, T1574, T1098",
  "content_type": "procedure",
  "category": "endpoint_security",
  "subcategory": "persistence_investigation",
  "tags": ["persistence", "autoruns", "services", "scheduled-tasks", "registry", "cron", "systemd", "incident-response", "forensics"],
  "severity_filter": ["medium", "high", "critical"],
  "incident_types": ["malware", "backdoor", "active_intrusion", "apt"],
  "ioc_types": ["hash", "ip", "domain"],
  "mitre_techniques": ["T1547.001", "T1547.006", "T1053.005", "T1546.003", "T1546.010", "T1546.012", "T1546.015", "T1543.003", "T1574.001", "T1053.003", "T1098.004"],
  "compliance_frameworks": ["NIST", "SOC2", "ISO27001"],
  "priority": 83
}
