# kb-content-output/articles/endpoint_security/windows-event-log-analysis.json
{
  "title": "Windows Security Event Log Analysis",
  "content": "## Purpose\n\nThis runbook provides detailed guidance for analyzing Windows Security Event Logs during incident investigation. Windows Event Logs are the foundational data source for detecting authentication abuse, privilege escalation, process execution, and persistence on Windows systems. This guide covers the most security-relevant Event IDs with practical analysis techniques and investigation queries.\n\n## Scope\n\nCovers analysis of Windows Security logs, System logs, PowerShell operational logs, and Sysmon logs. Applicable to incident response investigations, threat hunting, and forensic analysis on Windows 10/11 workstations and Windows Server 2019/2022.\n\n## Prerequisites\n\n- SIEM access with Windows Event Log data ingested\n- Direct endpoint access via EDR remote shell or RDP for live analysis\n- Event Viewer proficiency or PowerShell `Get-WinEvent` skills\n- Windows audit policies configured per CIS Benchmarks (see Endpoint Hardening Checklist)\n- Sysmon deployed for enhanced endpoint telemetry\n\n## Runbook Steps\n\n### Step 1: Authentication Event Analysis\n\n**Event ID 4624 - Successful Logon**\n\nThis is the most frequently analyzed event for tracking user access.\n\nKey fields to examine:\n- **Logon Type**: Determines HOW the user logged in (see Logon Type table below)\n- **Target User Name / Domain**: WHO logged in\n- **Source Network Address**: WHERE the logon originated from (blank for local logons)\n- **Logon Process / Authentication Package**: NTLM vs. Kerberos (NTLM in a Kerberos environment is suspicious)\n- **Elevated Token**: Indicates if the session has admin privileges\n\n| Logon Type | Name | Investigation Focus |\n|---|---|---|\n| 2 | Interactive | Local console login - expected for physical access |\n| 3 | Network | SMB, PsExec, WMI remote - check if source IP is authorized admin workstation |\n| 4 | Batch | Scheduled task - verify the task is legitimate |\n| 5 | Service | Service startup - correlate with Event 7045 for new services |\n| 7 | Unlock | Workstation unlock - useful for physical presence timeline |\n| 8 | NetworkCleartext | IIS Basic Auth - credentials sent in cleartext, review for misconfig |\n| 9 | NewCredentials | RunAs /netonly - potential pass-the-hash indicator |\n| 10 | RemoteInteractive | RDP - verify source is authorized and user is permitted |\n| 11 | CachedInteractive | Cached credentials - DC was unreachable, check network connectivity |\n\n**Investigation query example** (KQL-style):\n```\nSecurityEvent\n| where EventID == 4624\n| where LogonType in (3, 10)  // Network and RDP logons\n| where SourceNetworkAddress != \"\"\n| where TargetUserName !endswith \"$\"  // Exclude machine accounts\n| summarize count(), dcount(Computer) by TargetUserName, SourceNetworkAddress, LogonType\n| where count_ > 10 or dcount_Computer > 3\n| sort by count_ desc\n```\n\n**Event ID 4625 - Failed Logon**\n\nKey additional fields:\n- **Status / Sub Status codes**: Reveal the specific failure reason\n\n| Status Code | Sub Status | Meaning |\n|---|---|---|\n| 0xC000006D | 0xC000006A | Bad password (brute force indicator) |\n| 0xC000006D | 0xC0000064 | User does not exist (username enumeration) |\n| 0xC0000234 | - | Account locked out |\n| 0xC0000072 | - | Account disabled |\n| 0xC000006F | - | Logon outside authorized hours |\n| 0xC0000070 | - | Logon from unauthorized workstation |\n| 0xC0000071 | - | Expired password |\n| 0xC000015B | - | Logon type not granted (policy restriction) |\n\n**Event ID 4648 - Explicit Credential Logon**\n\nOccurs when a user provides credentials different from their current session (RunAs, mapping drives with different credentials). Investigate when:\n- A standard user provides domain admin credentials\n- The target server is a domain controller or sensitive system\n- The source and target accounts are in different trust domains\n\n### Step 2: Account and Group Management Events\n\n| Event ID | Description | Investigation Trigger |\n|---|---|---|\n| 4720 | User account created | Any creation outside HR onboarding process |\n| 4722 | User account enabled | Re-enabled accounts, especially if previously disabled for cause |\n| 4724 | Password reset attempt | Resets for accounts the analyst does not own |\n| 4725 | User account disabled | Verify authorized by IT/HR |\n| 4726 | User account deleted | Verify authorized; may indicate covering tracks |\n| 4728 | Member added to global security group | **Critical**: Domain Admins, Enterprise Admins additions |\n| 4732 | Member added to local security group | Local Administrators group additions |\n| 4756 | Member added to universal security group | Schema Admins, Enterprise Admins |\n| 4738 | User account changed | Review for attribute changes (SPN set for Kerberoasting setup) |\n| 4740 | Account locked out | Track source of lockout for brute force correlation |\n| 4767 | Account unlocked | Verify authorized |\n\n**Critical group monitoring query**:\n```\nSecurityEvent\n| where EventID in (4728, 4732, 4756)\n| where TargetUserName in (\"Domain Admins\", \"Enterprise Admins\", \"Schema Admins\", \"Administrators\", \"Account Operators\", \"Backup Operators\")\n| project TimeGenerated, EventID, SubjectUserName, MemberName, TargetUserName, Computer\n```\n\n### Step 3: Process Execution Analysis\n\n**Event ID 4688 - Process Creation**\n\nRequires \"Audit Process Creation\" policy and \"Include command line in process creation events\" enabled.\n\nKey fields:\n- **New Process Name**: The executable that was launched\n- **Creator Process Name**: The parent process (critical for process tree analysis)\n- **Command Line** (Process Command Line): The full command with arguments\n- **Token Elevation Type**: 1 = Full token (admin), 2 = Elevated, 3 = Limited\n\n**Suspicious process creation patterns**:\n```\n# Hunt for encoded PowerShell\nSecurityEvent\n| where EventID == 4688\n| where NewProcessName endswith \"powershell.exe\" or NewProcessName endswith \"pwsh.exe\"\n| where CommandLine contains_any (\"-enc\", \"-EncodedCommand\", \"FromBase64\", \"IEX\", \"DownloadString\")\n\n# Hunt for LOLBin abuse\nSecurityEvent\n| where EventID == 4688\n| where NewProcessName has_any (\"certutil.exe\", \"mshta.exe\", \"regsvr32.exe\", \"bitsadmin.exe\")\n| where CommandLine contains_any (\"http\", \"ftp\", \"urlcache\", \"/transfer\")\n\n# Hunt for credential dumping tools\nSecurityEvent\n| where EventID == 4688\n| where CommandLine contains_any (\"mimikatz\", \"sekurlsa\", \"lsadump\", \"procdump\", \"-ma lsass\", \"comsvcs.dll MiniDump\")\n```\n\n### Step 4: Service and Scheduled Task Events\n\n**Event ID 7045 (System Log) - New Service Installed**\n\n| Field | Normal | Suspicious |\n|---|---|---|\n| Service Name | Descriptive (e.g., \"Windows Update\") | Random characters (e.g., \"xyzabc\") |\n| Service File Name | `C:\\Windows\\System32\\` or `C:\\Program Files\\` | `C:\\Users\\`, `C:\\Temp\\`, `%APPDATA%` |\n| Service Type | kernel driver or user mode service | user mode service running cmd.exe |\n| Service Start Type | auto start or demand start | auto start for unknown service |\n\n**Event ID 4698 - Scheduled Task Created**\n\nReview the task XML embedded in the event data for:\n- Tasks running as SYSTEM from non-admin sources\n- Tasks executing from temporary or user-writable directories\n- Tasks with encoded commands or download cradles\n- Tasks created outside change management windows\n\n### Step 5: PowerShell Logging Analysis\n\n**PowerShell Script Block Logging** (Event ID 4104 in Microsoft-Windows-PowerShell/Operational):\n- Captures the actual script content, even if it was encoded on the command line\n- Search for known malicious cmdlets: `Invoke-Mimikatz`, `Invoke-Kerberoast`, `Get-GPPPassword`, `Invoke-BloodHound`\n- Search for AMSI bypass attempts: `AmsiUtils`, `amsiInitFailed`, `Reflection.Assembly`\n- Search for download functions: `Net.WebClient`, `Invoke-WebRequest`, `Start-BitsTransfer`, `System.Net.HttpWebRequest`\n\n**PowerShell Module Logging** (Event ID 4103):\n- Records the output of PowerShell commands\n- Useful for seeing what data was returned from discovery commands\n\n### Step 6: Sysmon Event Analysis\n\nSysmon provides superior telemetry when standard Windows events are insufficient.\n\n**Priority Sysmon events for investigation**:\n\n1. **Event ID 1 (Process Create)**: Includes SHA256 hash, parent process, and full command line. Superior to Event 4688.\n2. **Event ID 3 (Network Connection)**: Maps processes to their network connections. Essential for C2 detection.\n3. **Event ID 8 (CreateRemoteThread)**: Process injection detection. Alert on any instance where source and target are different processes.\n4. **Event ID 10 (Process Access)**: Credential theft detection. Filter for target image = `lsass.exe` with access masks `0x1010`, `0x1038`, `0x1FFFFF`.\n5. **Event ID 11 (File Create)**: Malware drop detection. Focus on executables created in temp directories.\n6. **Event ID 13 (Registry Value Set)**: Persistence detection. Monitor Run keys and service modifications.\n7. **Event ID 22 (DNS Query)**: Process-level DNS. Identify which process queried a suspicious domain.\n8. **Event ID 25 (Process Tampering)**: Direct detection of process hollowing and herpaderping.\n\n## Common Investigation Workflows\n\n**Workflow 1: Compromised Account Investigation**\n1. Event 4625 (failures) > Event 4624 (success) from same source = brute force success\n2. Event 4624 > Event 4672 (special privileges) = privileged access gained\n3. Event 4672 > Event 4728/4732 (group changes) = privilege escalation\n\n**Workflow 2: Lateral Movement Investigation**\n1. Event 4624 Type 3 from workstation to server = network logon\n2. Event 7045 on target = service-based execution (PsExec)\n3. Event 4688 on target = command execution\n\n**Workflow 3: Persistence Investigation**\n1. Event 7045 (new service) or 4698 (new task) = persistence established\n2. Event 4688 (process creation from service/task) = persistence executed\n3. Sysmon Event 13 (registry change) = registry-based persistence\n\n## Escalation Criteria\n\n- Event 1102 (Audit Log Cleared): P1 - Anti-forensics activity, immediate investigation\n- Event 4728/4732 adding unexpected member to Domain Admins: P1 - Privilege escalation\n- Event 4688 showing Mimikatz or credential dumping tools: P1 - Active credential theft\n- Event 7045 with suspicious service on domain controller: P1 - DC compromise\n- Sysmon Event 10 targeting lsass.exe from unknown process: P1 - Credential access attempt\n\n## References\n\n- Microsoft Windows Security Event Documentation\n- SANS Windows Security Event Log Quick Reference\n- Sysmon Documentation: https://docs.microsoft.com/en-us/sysinternals/downloads/sysmon\n- NSA Event Forwarding Guidance\n- JPCERT Windows Event Log Analysis Guide",
  "content_type": "runbook",
  "category": "endpoint_security",
  "subcategory": "log_analysis",
  "tags": ["windows-events", "event-ids", "security-log", "sysmon", "powershell-logging", "4624", "4625", "4688", "7045", "forensics"],
  "severity_filter": ["low", "medium", "high", "critical"],
  "incident_types": ["unauthorized_access", "credential_theft", "lateral_movement", "privilege_escalation", "malware"],
  "ioc_types": ["ip", "hash"],
  "mitre_techniques": ["T1078", "T1110", "T1003", "T1059.001", "T1021", "T1543.003", "T1053.005", "T1070.001", "T1547.001"],
  "compliance_frameworks": ["NIST", "SOC2", "ISO27001", "PCI-DSS", "HIPAA"],
  "priority": 87
}
